# packages/sysutils/sunxi-sys-utils/system.d/sunxi-bt-load.service
[Unit]
Description=bluetooth firmware loader
After=network.target

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=-/bin/sh -c ". /etc/profile; exec /bin/sh /usr/lib/openelec/bt_init"

[Install]
WantedBy=network.target
