# packages/sysutils/sunxi-sys-utils/system.d/sunxi-poweroff.service
[Unit]
Description=sunxi poweroff service
Before=systemd-poweroff.service
After=kodi-poweroff.service
DefaultDependencies=no

[Service]
Type=oneshot
Environment=HOME=/storage
ExecStart=-/bin/sh -c "echo 1 > /sys/class/leds/green_led/brightness"
StandardError=null
RemainAfterExit=yes

[Install]
WantedBy=poweroff.target
