# packages/mediacenter/xbmc/system.d/display-manager.service
[Unit]
Description=Fake Display Manager
Requires=multi-user.target
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/bin/true
RemainAfterExit=yes

[Install]
WantedBy=graphical.target
