# packages/mediacenter/xbmc/system.d/kodi-sources.service
[Unit]
Description=Kodi sources Setup
Before=kodi.service

[Service]
Type=oneshot
Environment=HOME=/storage
ExecStart=/usr/lib/kodi/kodi-sources
RemainAfterExit=yes

[Install]
WantedBy=kodi.service
