# packages/mediacenter/xbmc/system.d/kodi-rebrand.service
[Unit]
Description=Kodi config dir rebrand
DefaultDependencies=no
Before=local-fs.target shutdown.target
Conflicts=shutdown.target

[Service]
Type=oneshot
ExecStart=/usr/lib/kodi/kodi-rebrand
RemainAfterExit=yes

[Install]
WantedBy=local-fs.target
