# packages/mediacenter/xbmc/system.d/kodi-cleanlogs.service
[Unit]
Description=Kodi clean debug logs
ConditionKernelCommandLine=!debugging
ConditionPathExists=!/storage/.cache/debug.openelec
Before=kodi.service

[Service]
Type=oneshot
ExecStart=-/bin/sh -c 'rm -rf /storage/.kodi/userdata/addon_data/*/*.log /storage/.kodi/userdata/addon_data/*/log/*'
RemainAfterExit=yes

[Install]
WantedBy=kodi.service
