# packages/mediacenter/xbmc/system.d/xbmc.service
[Unit]
Description=XBMC Media Center
After=xorg.target graphical.target
Requires=graphical.target

[Service]
Environment=DISPLAY=:0.0
Environment=SDL_MOUSE_RELATIVE=0
Environment=HOME=/storage
EnvironmentFile=-/run/xbmc/xbmc.conf
ExecStart=/usr/lib/xbmc/xbmc.bin --standalone -fs $XBMC_ARGS
KillMode=process
Restart=always

[Install]
WantedBy=xbmc.target
