# packages/mediacenter/xbmc/system.d/xbmc-lcd-suspend.service
[Unit]
Description=LCD sleep hook
Before=sleep.target
StopWhenUnneeded=yes

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=-/usr/bin/xbmc-send --host=127.0.0.1 -a "LCD.Suspend"
ExecStop=-/usr/bin/xbmc-send --host=127.0.0.1 -a "LCD.Resume"

[Install]
WantedBy=sleep.target
