# packages/mediacenter/xbmc/system.d/xbmc-cleanlogs.service
[Unit]
Description=XBMC clean debug logs
ConditionKernelCommandLine=!debugging
Before=xbmc.service

[Service]
Type=oneshot
Environment=HOME=/storage
ExecStart=-/bin/sh -c rm -rf /storage/.xbmc/userdata/addon_data/*/*.log /storage/.xbmc/userdata/addon_data/*/log/*
RemainAfterExit=yes

[Install]
WantedBy=xbmc.service
