# packages/mediacenter/xbmc/system.d/xbmc-userfonts.service
[Unit]
Description=XBMC user fonts setup
Before=xbmc.service

[Service]
Type=oneshot
ExecStart=/usr/lib/xbmc/xbmc-userfonts
RemainAfterExit=yes

[Install]
WantedBy=xbmc.service
