# packages/mediacenter/xbmc/system.d/xbmc-waitonnetwork.service
[Unit]
Description=Wait on network
Requisite=connman.service
After=connman.service
Before=xorg-server.service framebuffer.service

ConditionPathExists=/storage/.cache/openelec/network_wait

[Service]
Type=oneshot
ExecStartPre=/bin/sh -c 'echo "waiting on Network to come online ..."
ExecStart=/usr/bin/cm-online --timeout=30
StandardOutput=tty

[Install]
WantedBy=xbmc.service
