# packages/mediacenter/xbmc/system.d/xbmc.target
[Unit]
Description=XBMC Mediacenter Interface
Requires=multi-user.target graphical.target
After=graphical.target
Conflicts=rescue.target
AllowIsolate=yes

[Install]
Alias=default.target
