# packages/linux-firmware/rockchip-firmware/system.d/hciattach-realtek.service
[Unit]
Description=Attach /dev/ttyS0 to BlueZ stack
Wants=bluetooth.service
Before=bluetooth.service
After=dev-ttyS0.device

[Service]
Type=simple
ExecStartPre=/usr/lib/libreelec/hciattach-realtek-config
ExecStartPre=/usr/sbin/rfkill unblock bluetooth
ExecStart=/usr/bin/rtk_hciattach -n -s 115200 ttyS0 rtk_h5
ExecStopPost=/usr/sbin/rfkill block bluetooth
