# packages/linux-firmware/brcmfmac_sdio-firmware-aml-all/system.d/brcmfmac_sdio-firmware-aml.service
[Unit]
Description=Broadcom sdio firmware update

[Service]
Type=simple
RemainAfterExit=yes
EnvironmentFile=/usr/share/bttty.conf
ExecStart=/usr/bin/hciattach -n -s 115200 ${BTTTY} bcm43xx 2000000
