# packages/linux-firmware/brcmfmac_sdio-firmware-aml-all/system.d/brcmfmac_sdio-firmware@.service
[Unit]
Description=Broadcom sdio firmware update for %I

[Service]
Type=simple
RemainAfterExit=yes
EnvironmentFile=/usr/share/bttty.conf
ExecStart=/usr/bin/brcm_patchram_plus --patchram /usr/lib/firmware/brcm/%I.hcd --baudrate 2000000 --use_baudrate_for_download ${BTTTY} --enable_hci --no2bytes --tosleep=50000
