# packages/sysutils/imx6-soc-fan/system.d/imx6-soc-fan-monitor.service
[Unit]
Description=i.MX6 SoC fan monitor service
ConditionPathExists=/var/run/system_type_matrix

[Service]
ExecStart=/usr/bin/imx6-soc-fan-monitor

[Install]
WantedBy=multi-user.target
