# packages/sysutils/imx6-soc-fan/system.d/imx6-soc-fan-monitor.service
[Unit]
Description=i.MX6 SoC fan monitor service

[Service]
ExecStart=/bin/sh -c "exec /bin/sh /usr/lib/openelec/imx6-soc-fan-monitor"
RestartPreventExitStatus=12
TimeoutStopSec=5
Restart=always
RestartSec=60
StartLimitInterval=0

[Install]
WantedBy=multi-user.target
