# packages/x11/other/fluxbox/system.d/windowmanager.service
[Unit]
Description=Fluxbox Window Manager
Requires=xorg.service
After=xorg.service
Before=xbmc.service

[Service]
Environment=DISPLAY=:0
Environment=HOME=/storage
ExecStart=/usr/bin/fluxbox
KillMode=process
Restart=always
RestartSec=1
StartLimitInterval=0
TimeoutStopSec=1s

[Install]
WantedBy=graphical.target
