# projects/L4T/packages/L4T/system.d/tmp-database.mount
[Unit]
Description=Database directory
Before=retroarch.service
After=storage.mount
After=systemd-tmpfiles-setup.service

[Mount]
What=/storage/database:/usr/share/libretro-database
Where=/tmp/database
Type=mergerfs
Options=defaults,allow_other,use_ino

[Install]
WantedBy=retroarch.target
