# projects/L4T/packages/openssl/system.d/openssl-config.service
[Unit]
Description=OpenSSL configuration service
DefaultDependencies=no
After=systemd-tmpfiles-setup.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/openssl-config

[Install]
WantedBy=sysinit.target
