# packages/sysutils/imx6-status-led/system.d/imx6-status-led.service
[Unit]
Description=i.MX6 status LED service
After=kodi.service

[Service]
Type=oneshot
ExecStart=-/bin/sh -c "/usr/bin/imx6-status-led on"
ExecStop=-/bin/sh -c "/usr/bin/imx6-status-led heartbeat"
RemainAfterExit=yes

[Install]
WantedBy=kodi.target
