# packages/sysutils/atvclient/system.d/atvclient.service
[Unit]
Description=ATV client remote support
After=syslog.target

[Service]
ExecStart=/usr/sbin/atvclient -m
KillMode=process

