# packages/sysutils/udisks/system.d/udisks-mount-disks.service
[Unit]
Description=Mount internal disks

[Service]
ExecStart=/usr/lib/udisks/udisks-mount-disks
Type=oneshot

[Install]
WantedBy=remote-fs.target
