# packages/sysutils/udisks/system.d/udisks.service
[Unit]
Description=Udisks server daemon
After=syslog.target polkit.service

[Service]
ExecStart=/usr/lib/udisks/udisks-daemon --no-debug
KillMode=process

[Install]
WantedBy=xbmc.target
