# packages/sysutils/udisks/system.d/udisks.service
[Unit]
Description=Disk Manager
Documentation=man:udisks(8)

[Service]
Type=dbus
BusName=org.freedesktop.UDisks
ExecStart=/usr/lib/udisks/udisks-daemon
