# projects/RPi/devices/RPi4-PiBoyDmg/packages/xpi_gamecon/system.d/piboy-dmg-control.service
[Unit]
Description=PiBoy DMG Control Service

[Service]
ExecStart=/usr/bin/piboy-dmg-control.sh
KillMode=killall piboy-dmg-control.sh 
TimeoutStopSec=10
Restart=always
RestartSec=2
StartLimitInterval=0

[Install]
WantedBy=multi-user.target
