# projects/RPi/devices/RPi4-PiBoyDmg/packages/xpi_gamecon/system.d/xpi_gamecon_reboot.service
[Unit]
Description=Shutdown xpi_gamecon
DefaultDependencies=no

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/bin/sh -c 'echo "129" > /sys/kernel/xpi_gamecon/flags; sleep 1'

[Install]
WantedBy=reboot.target
