# projects/RPi/devices/RPi4-PiBoyDmg/packages/xpi_gamecon/system.d/xpi_gamecon_shutdown.service
[Unit]
Description=Shutdown xpi_gamecon
DefaultDependencies=no
Before = local-fs.target

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/bin/sh -c 'echo "0" > /sys/kernel/xpi_gamecon/flags'
ExecStop=/usr/sbin/rmmod xpi_gamecon

[Install]
WantedBy=halt.target poweroff.target
