# projects/Rockchip/devices/OdroidGoAdvance/packages/odroidgoa-volkeys/system.d/odroidgoa-volkeys.service
[Unit]
Description=Odroid Go Advance volume keys sense service
Before=retroarch.target

[Service]
Environment=HOME=/storage
ExecStart=/usr/bin/spkeys-service odroidgo3-keys
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
