# projects/L4T/packages/systemd/system.d/systemd-timesyncd.service.d/depend-on-kernel-ip-config.conf
[Unit]
ConditionPathExists=/dev/.kernel_ipconfig
After=systemd-timesyncd-setup.service network-base.service
