# packages/sysutils/remote/atvclient/udev.d/98-appleir.rules
# must be started in a shell, otherwise it exits after sometime for unknown reasons
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="05ac", ATTR{idProduct}=="824[012]", RUN+="/bin/sh -c '/usr/sbin/atvclient -m &'"
