# projects/Switch/packages/switch-gpu-profile/system.d/switch-gpu-profile.service
[Unit]
Description=Switch GPU Profile
ConditionKernelCommandLine=!installer

[Service]
Environment=TERM=linux
Type=oneshot
ExecStart=/usr/bin/gpu-profile init
StandardError=syslog
RemainAfterExit=no

[Install]
WantedBy=graphical.target
