# projects/L4T/devices/Switch/packages/usb-gadget-scripts/udev.d/92-usb-gadgets.rules
SUBSYSTEM!="extcon", GOTO="gadget_end"
KERNEL!="extcon1", GOTO="gadget_end"
ATTRS{cable.1/state}=="1", TEST=="/usr/bin/usb-gadget.sh", RUN+="/usr/bin/usb-gadget.sh finalize"
ATTRS{cable.1/state}=="0", TEST=="/usr/bin/usb-gadget.sh", RUN+="/usr/bin/usb-gadget.sh disconnect"
LABEL="gadget_end"
