# packages/addons/service/webgrabplus/source/system.d/service.webgrabplus.service
[Unit]
Description=WebGrab+Plus
Documentation=http://www.webgrabplus.com/
After=network.target
Requires=network.target

[Service]
ExecStart=/bin/sh /storage/.kodi/addons/service.webgrabplus/bin/webgrabplus.start
Restart=always
RestartSec=21600

[Install]
WantedBy=kodi.target
