# packages/addons/service/webgrabplus/source/system.d/service.webgrabplus.service
[Unit]
Description=WebGrab+Plus
Documentation=http://www.webgrabplus.com/
After=network.target
Requires=network.target

[Service]
ExecStart=/bin/sh /storage/.kodi/addons/service.webgrabplus/bin/webgrabplus.start
ExecStop=/usr/bin/systemctl stop webgrabplus.timer
ExecStop=/usr/bin/systemctl stop webgrabplus.service
RemainAfterExit=yes
Type=oneshot

[Install]
WantedBy=kodi.target
