# packages/sysutils/udev/udev.d/80-clock.rules
ACTION=="add", SUBSYSTEM=="rtc", RUN+="/sbin/hwclock --hctosys --rtc=/dev/%k"
ACTION=="add", MAJOR==10, MINOR==135, RUN+="/sbin/hwclock --hctosys --rtc=/dev/%k"
