# packages/addons/service/tigervnc/source/system.d/service.tigervnc.service
[Unit]
Description=TigerVNC
After=graphical.target

[Service]
Environment=DISPLAY=:0.0
ExecStart=/bin/sh /storage/.kodi/addons/service.tigervnc/bin/tigervnc.start
TimeoutStopSec=1
Restart=always
RestartSec=2
StartLimitInterval=0

[Install]
WantedBy=kodi.target
