# packages/mediacenter/service.openelec.settings/system.d/backup-restore.service
[Unit]
Description=Restoring Backup
DefaultDependencies=false
After=show-version.service
Before=local-fs-pre.target shutdown.target
Conflicts=shutdown.target

[Service]
Type=oneshot
ExecStart=/usr/lib/openelec/backup-restore
StandardOutput=tty

[Install]
WantedBy=local-fs-pre.target

