# packages/mediacenter/service.openelec.settings/system.d/backup-restore.service
[Unit]
Description=Restoring Backup
DefaultDependencies=false
After=storage.mount show-version.service
Before=local-fs-pre.target swap.target

[Service]
Type=oneshot
ExecStart=/usr/lib/openelec/backup-restore
StandardOutput=tty

[Install]
WantedBy=basic.target
