# packages/addons/service/net-snmp/source/system.d/service.net-snmp.service
[Unit]
Description=Net-SNMP software bundle
After=network-online.target
Requires=network-online.target

[Service]
Type=forking
ExecStart=/bin/sh -c "exec sh /storage/.kodi/addons/service.net-snmp/bin/snmpd.start"
ExecStop=/bin/sh -c "exec sh /storage/.kodi/addons/service.net-snmp/bin/snmpd.stop"

[Install]
WantedBy=kodi.target
