# packages/network/iw/udev.d/60-iw-regdomain.rules
SUBSYSTEM=="ieee80211", ACTION=="add", RUN+="/usr/lib/iw/setregdomain"
