# projects/S802/packages/linux-firmware/brcmfmac_sdio-firmware-aml/system.d/brcmfmac_sdio-firmware@.service
[Unit]
Description=Broadcom sdio firmware update for %I
ConditionPathExists=/dev/ttyS2

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=/usr/bin/brcm_patchram_plus --patchram /usr/lib/firmware/brcm/%I.hcd --baudrate 3000000 --use_baudrate_for_download /dev/ttyS2 --enable_hci --no2bytes --tosleep=50000
