# projects/Generic/devices/x11/filesystem/usr/lib/systemd/system/intel-fullrange-resume.service
[Unit]
Description=Restore full range after suspend
Before=sleep.target
StopWhenUnneeded=yes

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/true
ExecStop=-/bin/sh "/usr/bin/intel-fullrange.sh"

[Install]
WantedBy=sleep.target
