# packages/addons/service/inadyn/source/system.d/service.system.inadyn.service
[Unit]
Description=inadyn
After=network-online.target
Requires=network-online.target

[Service]
Type=forking
ExecStart=/bin/sh /storage/.kodi/addons/service.system.inadyn/bin/inadyn-service
Restart=on-failure

[Install]
WantedBy=kodi.target
