# packages/addons/driver/sundtek-mediatv/source/system.d/driver.dvb.sundtek-mediatv.service
[Unit]
Description=Sundtek driver service
After=network.target

[Service]
Type=oneshot
ExecStart=-/bin/sh -c "exec sh /storage/.kodi/addons/driver.dvb.sundtek-mediatv/bin/sundtek-mediatv.start"
ExecStop=-/bin/sh -c "exec sh /storage/.kodi/addons/driver.dvb.sundtek-mediatv/bin/sundtek-mediatv.stop"
RemainAfterExit=yes

[Install]
WantedBy=kodi.target
