# packages/tools/xbox360-controllers-shutdown/system.d/xbox360-controllers-shutdown.service
[Unit]
Description=xbox360-controllers-shutdown
DefaultDependencies=no
Conflicts=reboot.target
Before=poweroff.target halt.target shutdown.target
Requires=poweroff.target shutdown.target

[Service]
Type=oneshot
ExecStart=$INSTALL/usr/bin/xbox360-controllers-shutdown
RemainAfterExit=yes

[Install]
WantedBy=shutdown.target
