# packages/addons/service/docker/source/examples/var-lib-docker.mount
[Unit]
Description=Docker loopback mount
Requires=local-fs.target
Before=service.system.docker.service

[Mount]
What=/storage/btrfs-loop.img
Where=/var/lib/docker

Options=loop,compress=lzo
Type=btrfs

[Install]
WantedBy=multi-user.target
