# packages/linux-firmware/amlogic/qca9377-firmware-aml/system.d/qca9377-firmware-aml.service
[Unit]
Description=QCA Bluetooth firmware service
ConditionPathExists=/dev/ttyS1

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=/usr/bin/hciattach -n -s 115200 /dev/ttyS1 qca 2000000
