# packages/addons/driver/steamcontroller/source/system.d/driver.steamcontroller.service
[Unit]
Description=Steam Controller driver

[Service]
Type=forking
ExecStart=/usr/bin/python /storage/.kodi/addons/driver.steamcontroller/bin/sc-xbox.py start
ExecReload=/usr/bin/python /storage/.kodi/addons/driver.steamcontroller/bin/sc-xbox.py restart
ExecStop=/usr/bin/python /storage/.kodi/addons/driver.steamcontroller/bin/sc-xbox.py stop
PIDFile=/run/steamcontroller.pid
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
