# projects/L4T/packages/lirc/system.d/lircd.socket
[Socket]
ListenStream=/run/lirc/lircd.socket

[Install]
WantedBy=sockets.target
Also=lircd.service
