# packages/linux/system.d/cpufreq-threshold.service
[Unit]
Description=Set ondemand threshold
DefaultDependencies=false

[Service]
Type=oneshot
ExecStart=-/bin/sh -c 'echo 50 > /sys/devices/system/cpu/cpufreq/ondemand/up_threshold'
RemainAfterExit=yes

[Install]
WantedBy=basic.target
