# packages/linux/system.d/cpupower-start.service
[Unit]
Description=Set cpupower to performance
DefaultDependencies=false

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'for i in /sys/devices/system/cpu/*/cpufreq/scaling_governor; do printf "performance" > $i; done'
RemainAfterExit=yes

[Install]
WantedBy=basic.target
