# packages/linux/system.d/cpupower-end.timer
[Unit]
Description=Set cpupower to ondemand

[Timer]
OnBootSec=40sec

[Install]
WantedBy=sysinit.target
