# packages/addons/service/emby4/source/system.d/service.emby4.service
[Unit]
Description=Emby 4 - streaming server
Documentation=http://emby.media
Wants=multi-user.target
After=multi-user.target
Conflicts=service.emby.service

[Service]
ExecStart=/bin/sh /storage/.kodi/addons/service.emby4/bin/emby4.start
Restart=always

[Install]
WantedBy=multi-user.target
