# packages/sysutils/open-iscsi/system.d/iscsi-initiator.service
[Unit]
Description=Open-iSCSI initiator (i.e. client) service
After=syslog.target

[Service]
ExecStart=/usr/sbin/iscsid -f

[Install]
WantedBy=multi-user.target
