# packages/addons/service/ttyd/source/system.d/service.ttyd.service
[Unit]
Description=ttyd
After=graphical.target

[Service]
ExecStart=/bin/sh -c "exec sh /storage/.kodi/addons/service.ttyd/bin/ttyd.start"
TimeoutStopSec=1
Restart=always
RestartSec=2
StartLimitInterval=0
Type=forking

[Install]
WantedBy=kodi.target
