# packages/sysutils/remote/eventlircd/system.d/eventlircd.service
[Unit]
Description=Eventlirc server daemon

[Service]
ExecStart=/usr/sbin/eventlircd -f --evmap=/etc/eventlircd.d --socket=/run/lirc/lircd --release=_UP
KillMode=process
TimeoutStopSec=1s

[Install]
WantedBy=multi-user.target

