# packages/sysutils/remote/eventlircd/system.d/eventlircd.service
[Unit]
Description=Eventlirc server daemon

[Service]
EnvironmentFile=-/run/openelec/debug/eventlircd.conf
ExecStart=/usr/sbin/eventlircd -f --evmap=/etc/eventlircd.d --socket=/var/run/lirc/lircd --release=_UP $EVENTLIRCD_DEBUG
KillMode=process

[Install]
WantedBy=multi-user.target

