# packages/sysutils/remote/eventlircd/system.d/eventlircd.service
[Unit]
Description=Eventlirc server daemon
After=syslog.target

[Service]
# EnvironmentFile=/etc/sysconfig/eventlircd
ExecStart=/usr/sbin/eventlircd -f --evmap=/etc/eventlircd.d --socket=/var/run/lirc/lircd --release=_UP
# ExecReload=/bin/kill -HUP $MAINPID
KillMode=process

[Install]
WantedBy=multi-user.target

