# packages/sysutils/remote/eventlircd/system.d/eventlircd.service
[Unit]
Description=Eventlirc server daemon
After=syslog.target

[Service]
EnvironmentFile=/etc/sysconf.d/eventlircd.conf
ExecStart=/usr/sbin/eventlircd -f $EVENTLIRCD_ARGS
KillMode=process

[Install]
WantedBy=multi-user.target

