# packages/sysutils/remote/lirc/system.d/lircd@.service
[Unit]
Description=Lirc with %I

[Service]
Type=oneshot
ExecStart=/usr/lib/openelec/lircd_helper add %I
ExecStop=/usr/lib/openelec/lircd_helper remove %I
TimeoutStopSec=1
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
