# packages/x11/other/fontconfig/system.d/fc-cache.service
[Unit]
Description=Fontconfig caching
Before=graphical.target

[Service]
Type=oneshot
ExecStart=/usr/bin/fc-cache /usr/share/fonts/ /usr/share/xbmc/media/Fonts/
RemainAfterExit=yes

[Install]
WantedBy=graphical.target
