# packages/addons/service/mariadb/source/system.d/service.mariadb.service
[Unit]
Description=MariaDB server
After=graphical.target network-online.service
Before=kodi.service

[Service]
Type=forking
PIDFile=/var/run/mysqld/mysqld.pid
ExecStart=/bin/sh -c "exec sh /storage/.kodi/addons/service.mariadb/bin/mariadb.start"
ExecStop=/bin/sh -c "exec sh /storage/.kodi/addons/service.mariadb/bin/mariadb.stop"
TimeoutStartSec=60
TimeoutStopSec=60
Restart=always
RestartSec=10
StartLimitInterval=0

[Install]
WantedBy=kodi.service
