# projects/Amlogic_Legacy/filesystem/usr/lib/systemd/system/fstrim.service
[Unit]
Description=Apply TRIM on all trimmable mounts

[Service]
Type=simple
ExecStart=-/usr/sbin/fstrim -a -v

[Install]
WantedBy=basic.target
