# packages/linux-firmware/wandboard-firmware/system.d/hciattach@.service
[Unit]
Description=Attach /dev/ttymxc2 to BlueZ stack
ConditionPathExists=/dev/ttymxc2

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=/usr/bin/hciattach /dev/ttymxc2 %I
