# projects/Switch/packages/busybox/system.d/show-version.service
[Unit]
Description=Show Version
DefaultDependencies=false
ConditionKernelCommandLine=!morequiet

[Service]
Type=oneshot
ExecStart=/usr/bin/lsb_release
RemainAfterExit=yes
StandardOutput=tty

[Install]
WantedBy=basic.target
