# packages/sysutils/remote/lirc/config/82_ircore.rules
#
# udev rules to start lirc on module load
#
# contributed by Stephan Raue <stephan@openelec.tv>
#

ACTION!=="add|change", GOTO="remote_event_grp_end"
SUBSYSTEM!="input", GOTO="remote_event_grp_end"

KERNEL=="event[0-9]*", SUBSYSTEMS=="rc", RUN+="/usr/sbin/lircd --driver=devinput --device=/dev/input/%k /etc/lirc/lircd.conf.devinput"
KERNEL=="event[0-9]*", SUBSYSTEMS=="input", ENV{.INPUT_CLASS}=="ir", RUN+="/usr/sbin/lircd --driver=devinput --device=/dev/input/%k /etc/lirc/lircd.conf.devinput"

LABEL="remote_event_grp_end"
