# projects/Ayn/devices/Odin/packages/retroarch/system.d/tmp-system.mount
[Unit]
Description=RetroArch system directory
Before=retroarch.service umount.target
Conflicts=umount.target
After=systemd-tmpfiles-setup.service
DefaultDependencies=no

[Mount]
What=/storage/system:/usr/share/retroarch-system
Where=/tmp/system
Type=mergerfs
Options=defaults,allow_other,use_ino

[Install]
WantedBy=retroarch.target
