# projects/Ayn/devices/Odin/packages/retroarch/system.d/tmp-cores.mount
[Unit]
Description=Cores directory
Before=retroarch.service umount.target
Conflicts=umount.target
After=systemd-tmpfiles-setup.service
DefaultDependencies=no

[Mount]
What=/storage/cores:/usr/lib/libretro
Where=/tmp/cores
Type=mergerfs
Options=defaults,allow_other,use_ino

[Install]
WantedBy=retroarch.target
