# projects/Ayn/devices/Odin/packages/retroarch/system.d/tmp-shaders.mount
[Unit]
Description=Shaders directory RetroArch
Before=retroarch.service umount.target
Conflicts=umount.target
After=systemd-tmpfiles-setup.service
DefaultDependencies=no

[Mount]
What=/storage/shaders:/usr/share/retroarch/shaders
Where=/tmp/shaders
Type=mergerfs
Options=defaults,allow_other,use_ino

[Install]
WantedBy=retroarch.target
