# packages/lakka/lakka_tools/wii-u-gc-adapter/system.d/wii-u-gc-adapter.service
[Unit]
Description=Wii U GC Adapter

[Service]
TimeoutStartSec=infinity
ExecStart=/usr/bin/wii-u-gc-adapter

[Install]
WantedBy=retroarch.target

