# packages/sysutils/irqbalanced/system.d/irqbalance.service
[Unit]
Description=irqbalance
After=proc.mount

[Service]
ExecStart=/usr/sbin/irqbalance
Type=forking

[Install]
WantedBy=sysinit.target
