# packages/network/avahi/system.d/avahi-monitor.path
[Unit]
Description=Avahi configfile monitor

[Path]
PathModified=/storage/.cache/services/avahi.conf
PathModified=/run/openelec/debug/avahi.conf

[Install]
WantedBy=multi-user.target
