# packages/network/avahi/system.d/avahi-defaults.service
[Unit]
Description=Avahi defaults
After=local-fs.target

ConditionPathExists=!/storage/.cache/services/avahi.conf
ConditionPathExists=!/storage/.cache/services/avahi.disabled

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'cp /usr/share/services/avahi.conf /storage/.cache/services/'
RemainAfterExit=yes
