# packages/network/avahi/system.d/avahi-daemon.service
[Unit]
Description=Avahi Daemon
After=syslog.target debug.service network.target

ConditionPathExists=/storage/.cache/services/avahi.conf

[Service]
EnvironmentFile=-/run/openelec/debug.conf
ExecStart=/usr/sbin/avahi-daemon -s $AVAHI_DEBUG
ExecReload=/usr/sbin/avahi-daemon -r

[Install]
WantedBy=multi-user.target
