# packages/network/avahi/system.d/avahi-monitor.service
[Unit]
Description=Avahi configfile monitor
After=avahi-daemon.service

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart avahi-daemon.service

[Install]
WantedBy=multi-user.target
