# packages/network/avahi/system.d/avahi-daemon.service
[Unit]
Description=Avahi Daemon
After=syslog.target avahi-defaults.service
Requires=avahi-defaults.service

ConditionPathExists=/storage/.cache/services/avahi.conf

[Service]
EnvironmentFile=-/run/openelec/debug/avahi.conf
ExecStart=/usr/sbin/avahi-daemon -s $AVAHI_DEBUG
ExecReload=/usr/sbin/avahi-daemon -r
TimeoutStopSec=1s

[Install]
WantedBy=multi-user.target
