# packages/addons/service/tvheadend42/source/system.d/service.tvheadend42.service
[Unit]
Description=TVHeadend42 Service
After=network-online.service

[Service]
ExecStart=/bin/sh -c "exec sh /storage/.kodi/addons/service.tvheadend42/bin/tvheadend42.start"
TimeoutStopSec=2
Restart=always
RestartSec=2
StartLimitInterval=0

[Install]
WantedBy=kodi.target
