# projects/Gamegirl/packages/gamegirl-screen/system.d/gamegirl-screen.service
[Unit]
Description=Gamegirl screen
After=graphical.target
Requires=graphical.target

[Service]
Type=oneshot
ExecStart=/usr/bin/enable_rgb565

[Install]
WantedBy=retroarch.target
