# packages/network/samba/system.d/samba-config.service
[Unit]
Description=Samba Configuration
DefaultDependencies=false
After=basic.target

[Service]
Type=oneshot
ExecStart=/usr/lib/samba/samba-config
StartLimitInterval=0
