# packages/network/samba/system.d/samba-config.service
[Unit]
Description=Samba Configuration
DefaultDependencies=false

[Service]
Type=oneshot
ExecStart=/usr/lib/samba/samba-config
StartLimitInterval=0

[Install]
WantedBy=sysinit.target
