# packages/lakka/lakka_tools/gpicase2_safeshutdown/system.d/gpicase2-safeshutdown.service
[Unit]
Description=GPICase 2 safe shutdown

[Service]
ExecStart=/usr/bin/safeshutdown_gpi2.py
TimeoutStopSec=5
Restart=always
RestartSec=5
StartLimitInterval=0

[Install]
WantedBy=multi-user.target
