# packages/lakka/Switch/system.d/mount-bluez-fs.service
[Unit]
Description=mount bluez filesystem
Requires=bluetooth.target
Before=bluetooth.service

[Service]
Type=oneshot
ExecStart=/usr/bin/mount-bluez-fs

[Install]
WantedBy=bluetooth.target
