# packages/lakka/Switch/system.d/switch-wifi-fix.service
[Unit]
Description=Switch Wireless Reboot Fix
After=connman.service
ConditionKernelCommandLine=!installer

[Service]
Environment=TERM=linux
Type=oneshot
ExecStart=/usr/bin/switch-wifi-fix
StandardError=syslog
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
