# packages/lakka/Switch/system.d/switch-set-mac-address.service
[Unit]
Description=set Switch mac address
Wants=network-pre.target
Before=network-pre.target

[Service]
Type=oneshot
ExecStart=/usr/bin/switch-set-mac-address

[Install]
WantedBy=multi-user.target
