# packages/lakka/Switch/system.d/tmp-database.mount
[Unit]
Description=Database directory
Before=retroarch.target
After=storage.mount
After=systemd-tmpfiles-setup.service

[Mount]
What=/usr/share/libretro-database:/storage/database
Where=/tmp/database
Type=mergerfs
Options=defaults,allow_other,use_ino

[Install]
WantedBy=retroarch.target
