# packages/lakka/Switch/system.d/tmp-shaders.mount
[Unit]
Description=Shaders directory RetroArch
Before=retroarch.target
After=storage.mount
After=systemd-tmpfiles-setup.service

[Mount]
What=/storage/shaders:/usr/share/common-shaders
Where=/tmp/shaders
Type=mergerfs
Options=defaults,allow_other,use_ino

[Install]
WantedBy=retroarch.target
