# packages/lakka/Switch/system.d/tmp-joypads.mount
[Unit]
Description=Joypad configs directory
Before=retroarch.service
After=storage.mount
After=systemd-tmpfiles-setup.service

[Mount]
What=/storage/joypads:/etc/retroarch-joypad-autoconfig
Where=/tmp/joypads
Type=mergerfs
Options=defaults,allow_other,use_ino

[Install]
WantedBy=retroarch.target
