# packages/lakka/Switch/system.d/tmp-cores.mount
[Unit]
Description=Cores directory
Before=retroarch.target
After=storage.mount
After=systemd-tmpfiles-setup.service

[Mount]
What=/storage/cores:/usr/lib/libretro
Where=/tmp/cores
Type=mergerfs
Options=defaults,allow_other,use_ino

[Install]
WantedBy=retroarch.target
