# packages/lakka/Switch/system.d/tmp-overlays.mount
[Unit]
Description=Overlays directory
Before=retroarch.target
After=storage.mount
After=systemd-tmpfiles-setup.service

[Mount]
What=/usr/share/retroarch-overlays:/storage/overlays
Where=/tmp/overlays
Type=mergerfs
Options=defaults,allow_other,use_ino

[Install]
WantedBy=retroarch.target
