# packages/lakka/Switch/system.d/xorg-configure-tegra.service
[Unit]
Description=configure Xorg Server for tegra

ConditionKernelCommandLine=!installer

[Service]
Type=oneshot
ExecStart=/usr/lib/xorg/xorg-configure tegra
RemainAfterExit=yes

[Install]
WantedBy=graphical.target
