# packages/addons/service/oscam/source/system.d/service.softcam.oscam.service
[Unit]
Description=OSCam
After=network-online.service
Requires=network-online.service

[Service]
ExecStart=/bin/sh -c "exec sh /storage/.kodi/addons/service.softcam.oscam/bin/oscam.start"
TimeoutStopSec=1
Restart=always
RestartSec=2
StartLimitInterval=0

[Install]
WantedBy=kodi.target
