# projects/Switch/packages/switch-cpu-profile/system.d/switch-cpu-profile.service
[Unit]
Description=Switch CPU Profile
ConditionKernelCommandLine=!installer

[Service]
Environment=TERM=linux
Type=oneshot
ExecStart=/usr/bin/cpu-profile init
StandardError=syslog
RemainAfterExit=no

[Install]
WantedBy=graphical.target
