# packages/network/connman/system.d/connman-vpn.service
[Unit]
Description=ConnMan VPN service
After=syslog.target

[Service]
Type=dbus
BusName=net.connman.vpn
EnvironmentFile=-/run/openelec/debug/connman.conf
ExecStart=/usr/sbin/connman-vpnd -n $CONNMANVPN_DEBUG
StandardOutput=null

[Install]
WantedBy=multi-user.target
