# packages/network/connman/system.d/connman-vpn.service
[Unit]
Description=ConnMan VPN service
Before=network.target
Wants=network.target

[Service]
Type=dbus
BusName=net.connman.vpn
ExecStart=/usr/sbin/connman-vpnd -n
StandardOutput=null

[Install]
WantedBy=multi-user.target
