# packages/network/connman/system.d/connman.service
[Unit]
Description=Connection service
After=debugconfig.service
Before=network.target
Wants=network.target

[Service]
Type=dbus
BusName=net.connman
Restart=on-failure
EnvironmentFile=-/run/openelec/debug/connman.conf
ExecStartPre=/usr/lib/connman/connman-setup
ExecStart=/usr/sbin/connmand -nr $CONNMAN_DEBUG
StandardOutput=null

[Install]
WantedBy=multi-user.target
