# packages/network/connman/system.d/connman.service
[Unit]
Description=Connection service
Before=network.target
Wants=network.target

[Service]
Type=dbus
BusName=net.connman
Restart=on-failure
ExecStart=/usr/sbin/connmand -n
StandardOutput=null

[Install]
WantedBy=multi-user.target
