# packages/network/connman/system.d/loopback.service
[Unit]
Description=Setup Loopback interface
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/sbin/ifconfig lo 127.0.0.1 netmask 255.0.0.0 up
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
