# packages/network/connman/system.d/loopback.service
[Unit]
Description=Setup Loopback interface
DefaultDependencies=no

[Service]
Type=oneshot
ExecStartPre=/sbin/ifconfig lo 127.0.0.1 netmask 255.0.0.0
ExecStart=/sbin/ifconfig lo up
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
