# packages/network/connman/system.d/network-online.service
[Unit]
Description=Wait for network to be configured by ConnMan
Requisite=connman.service
After=connman.service
Before=network-online.target
DefaultDependencies=no
Conflicts=shutdown.target

[Service]
Type=oneshot
ExecStartPre=/bin/sh -c 'echo "waiting on Network to come online ..."'
ExecStart=/usr/sbin/connmand-wait-online --timeout=30
StandardOutput=tty
RemainAfterExit=yes

[Install]
WantedBy=network-online.target
