# packages/network/connman/system.d/network-online.service
[Unit]
Description=Connection Manager Wait Online
Requisite=connman.service
After=connman.service

[Service]
Type=oneshot
ExecStartPre=/bin/sh -c 'echo "waiting on Network to come online ..."'
ExecStart=/usr/bin/cm-online 30
StandardOutput=tty

[Install]
WantedBy=multi-user.target
