# packages/addons/service/lcdd/source/system.d/service.lcdd.service
[Unit]
Description=LCDProc
After=multi-user.target

[Service]
ExecStart=/bin/sh /storage/.kodi/addons/service.lcdd/bin/lcdd.start
TimeoutStopSec=1s
Restart=on-failure
RestartSec=5
StartLimitInterval=30
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
