# packages/graphics/bcm2835-driver/system.d/bcm2835-threshold.service
[Unit]
Description=Set onedemand threshold
After=syslog.target

ConditionPathExists=/sys/devices/system/cpu/cpufreq/ondemand/up_threshold

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'echo 50 > /sys/devices/system/cpu/cpufreq/ondemand/up_threshold'

[Install]
WantedBy=multi-user.target

