# packages/graphics/bcm2835-driver/system.d/framebuffer.service
[Unit]
Description=Framebuffer setup
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/lib/bcm2835/setup-framebuffer
RemainAfterExit=yes

[Install]
WantedBy=graphical.target
