# packages/graphics/bcm2835-driver/system.d/vcfiled.service
[Unit]
Description=vcfiled
After=syslog.target

[Service]
ExecStart=/usr/sbin/vcfiled -f
KillMode=process

[Install]
WantedBy=multi-user.target

