# packages/network/samba/system.d.opt/nmbd.service
[Unit]
Description=Samba NMB Daemon
After=network.target smbd.service

[Service]
Type=forking
PIDFile=/var/run/nmbd-smb.conf.pid
ExecStart=/usr/bin/nmbd --configfile=/run/samba/smb.conf
ExecReload=/bin/kill -HUP $MAINPID
TimeoutStopSec=1s
Restart=always
RestartSec=2
StartLimitInterval=0

[Install]
WantedBy=multi-user.target
