# packages/network/samba/system.d.opt/samba-config.service
[Unit]
Description=Samba configfile writer
After=local-fs.target

ConditionPathExists=/storage/.cache/services/samba.conf

[Service]
Type=oneshot
ExecStart=/usr/lib/samba/samba-config

