# packages/network/samba/system.d.opt/samba-defaults.service
[Unit]
Description=Samba defaults
After=local-fs.target

ConditionPathExists=!/storage/.cache/services/samba.conf
ConditionPathExists=!/storage/.cache/services/samba.disabled

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'cp /usr/share/services/samba.conf /storage/.cache/services/'
RemainAfterExit=yes
