# packages/network/samba/system.d.opt/samba-config.service
[Unit]
Description=Samba configfile writer
After=samba-defaults.service
Requires=samba-defaults.service

ConditionPathExists=/storage/.cache/services/samba.conf

[Service]
Type=oneshot
ExecStart=/usr/lib/samba/samba-config
StartLimitInterval=0
