# packages/network/samba/system.d.opt/smbd-monitor.path
[Unit]
Description=Samba configfile monitor

[Path]
PathModified=/storage/.cache/services/samba.conf
PathModified=/storage/.config/samba.conf

[Install]
WantedBy=multi-user.target
