# packages/network/samba/system.d.opt/smbd.service
[Unit]
Description=Samba SMB Daemon
After=network.target samba-config.service
ConditionPathExists=/storage/.cache/services/samba.conf
Requires=samba-config.service

[Service]
Type=forking
PIDFile=/run/samba/smbd.pid
LimitNOFILE=16384
ExecStartPre=/usr/lib/samba/smbd-config
ExecStart=/usr/sbin/smbd
ExecReload=/bin/kill -HUP $MAINPID
TimeoutStopSec=1s
Restart=always
RestartSec=2
StartLimitInterval=0

[Install]
WantedBy=multi-user.target
