# packages/network/samba/system.d.opt/smbd-monitor.service
[Unit]
Description=Samba configfile monitor
After=network.target

[Service]
Type=oneshot
ExecStartPre=/bin/systemctl restart samba-config.service
ExecStartPre=/bin/systemctl restart nmbd.service
ExecStart=/bin/systemctl restart smbd.service

[Install]
WantedBy=multi-user.target
