# packages/sysutils/v4l-utils/udev.d/70-input-repeat.rules
ACTION=="add", SUBSYSTEM=="input", KERNEL=="event[0-9]*", GOTO="start"
GOTO="end"
LABEL="start"

# don't change eventlircd device settings
ATTRS{name}=="eventlircd", GOTO="end"

# don't change settings of rc remotes
SUBSYSTEMS=="rc", GOTO="end"

# set default repeat delay to 500ms like rc remotes
RUN+="/usr/bin/ir-keytable --delay=500 --device=$devnode"

# make lircd-uinput use same repeat period like keyboards in kodi
ATTRS{name}=="lircd-uinput", \
  RUN+="/usr/bin/ir-keytable --period=80 --device=$devnode"

LABEL="end"
