# packages/mediacenter/kodi/system.d/kodi.target
[Unit]
Description=Kodi Mediacenter Interface
Requires=multi-user.target graphical.target network-online.target
After=network-online.target graphical.target
Wants=network-online.target
Conflicts=rescue.target
AllowIsolate=yes

[Install]
Alias=default.target
