# packages/sysutils/busybox/system.d.opt/cron.service
[Unit]
Description=Cron daemon
After=syslog.target cron-defaults.service
Requires=cron-defaults.service

ConditionPathExists=/storage/.cache/services/crond.conf

[Service]
ExecStartPre=/bin/sh -c 'mkdir -p /storage/.cache/cron/crontabs'
ExecStart=/sbin/crond -f -S
KillMode=process
TimeoutStopSec=1s

[Install]
WantedBy=multi-user.target

