# packages/sysutils/busybox/system.d.opt/cron.service
[Unit]
Description=Cron daemon

ConditionPathExists=/storage/.cache/services/crond.conf

[Service]
ExecStartPre=/bin/sh -c 'mkdir -p /storage/.cache/cron/crontabs'
ExecStart=/sbin/crond -f -S
KillMode=process

[Install]
WantedBy=multi-user.target

