# packages/sysutils/udevil/system.d/udevil-mount@.service
[Unit]
Description=Udevil mount service

[Service]
Type=oneshot
ExecStart=-/usr/bin/udevil --mount %I
ExecStop=-/usr/bin/udevil --umount %I
RemainAfterExit=yes
