# packages/libretro/retroarch/system.d/localap.service
[Unit]
Description=Local AP via Connman
ConditionKernelCommandLine=|localap
ConditionPathExists=|/storage/.cache/services/localap.conf

[Service]
EnvironmentFile=-/storage/.cache/services/localap.conf
#ExecStart=/usr/bin/sshd -D $SSH_ARGS
ExecStart=/usr/bin/connmanctl tether wifi on $APNAME $PASSWORD
ExecStop=/usr/bin/connmanctl tether wifi off
TimeoutStopSec=1s
StartLimitInterval=0

[Install]
WantedBy=multi-user.target
