# packages/libretro/retroarch/system.d/tmp-database.mount
[Unit]
Description=Database directory
Before=retroarch.service
After=storage.mount
After=systemd-tmpfiles-setup.service

[Mount]
What=none
Where=/tmp/database
Type=overlay
Options=lowerdir=/usr/share/libretro-database,upperdir=/storage/database,workdir=/storage/.tmp/database-workdir

[Install]
WantedBy=retroarch.target
