# packages/libretro/retroarch/system.d/tmp-assets.mount
[Unit]
Description=Assets directory
Before=retroarch.service
After=storage.mount
After=systemd-tmpfiles-setup.service

[Mount]
What=none
Where=/tmp/assets
Type=overlay
Options=lowerdir=/usr/share/retroarch-assets,upperdir=/storage/assets,workdir=/storage/.tmp/assets-workdir

[Install]
WantedBy=retroarch.target
