# packages/libretro/retroarch/system.d/tmp-joypads.mount
[Unit]
Description=Joypad configs directory
Before=retroarch.service
After=storage.mount
After=systemd-tmpfiles-setup.service

[Mount]
What=none
Where=/tmp/joypads
Type=overlay
Options=lowerdir=/etc/retroarch-joypad-autoconfig,upperdir=/storage/joypads,workdir=/storage/.tmp/joypads-workdir

[Install]
WantedBy=retroarch.target
