# packages/libretro/retroarch/system.d/tmp-system.mount
[Unit]
Description=RetroArch system directory
Before=retroarch.service
After=systemd-tmpfiles-setup.service

[Mount]
What=none
Where=/tmp/system
Type=overlay
Options=lowerdir=/usr/share/retroarch-system,upperdir=/storage/system,workdir=/storage/.tmp/system-workdir

[Install]
WantedBy=retroarch.target
