# packages/libretro/retroarch/system.d/tmp-shaders.mount
[Unit]
Description=Shaders directory RetroArch
Before=retroarch.service
After=storage.mount
After=systemd-tmpfiles-setup.service

[Mount]
What=none
Where=/tmp/shaders
Type=overlay
Options=lowerdir=/usr/share/common-shaders,upperdir=/storage/shaders,workdir=/storage/.tmp/shaders-workdir

[Install]
WantedBy=retroarch.target
