# packages/libretro/retroarch/system.d/tmp-cores.mount
[Unit]
Description=Cores directory
Before=retroarch.service
After=storage.mount
After=systemd-tmpfiles-setup.service

[Mount]
What=none
Where=/tmp/cores
Type=overlay
Options=lowerdir=/usr/lib/libretro,upperdir=/storage/cores,workdir=/storage/.tmp/cores-workdir

[Install]
WantedBy=retroarch.target
