# projects/L4T/packages/tegra-bsp/assets/alsa-fix.service
[Unit]
Description=Fix alsa settings on Switch

[Service]
ExecStart=amixer sset "I2S1 Mux",0 ADMAIF1

[Install]
WantedBy=multi-user.target
