# projects/Rockchip/devices/OdroidGoAdvance/packages/odroidgo2-utils/system.d/odroidgoa-utils.service
[Unit]
Description=Odroid Go Advance headphone sense
ConditionPathExists=/dev/input/event1
Before=retroarch.target

[Service]
Environment=HOME=/storage
ExecStart=/usr/bin/headphone_sense.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
