# projects/L4T/devices/Switch/packages/switch-joycon-bluetooth-dock-configs/system.d/fix-permissions.service
[Unit]
Description=Fix permissions on sysfs nodes.
Before=graphical.target kodi.service retroarch.service
After=multi-user.target
After=xorg.service

[Service]
Environment=TERM=linux
Type=oneshot
ExecStart=/usr/bin/fix-sysfs-permissions.sh
StandardError=journal
RemainAfterExit=no

[Install]
WantedBy=graphical.target
