# projects/L4T/devices/Switch/packages/switch-joycon-bluetooth-dock-configs/system.d/fix-permissions.service
[Unit]
Description=Fix permissions on CPU/GPU freq policy
Before=mediacenter.target

[Service]
Environment=TERM=linux
Type=oneshot
ExecStart=/usr/bin/fix-freq-permissions.sh
StandardError=syslog
RemainAfterExit=no

[Install]
After=xorg.service
WantedBy=multi-user.target
