# projects/L4T/devices/Switch/packages/switch-joycon-bluetooth-dock-configs/system.d/pair-joycon.service
[Unit]
Description=Create Joycon Configs.
After=var-bluetoothconfig.mount
Before=bluetooth.service

[Service]
Environment=TERM=linux
Type=oneshot
ExecStart=/usr/bin/pair-joycon.sh
StandardError=syslog
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
