# projects/Rockchip/packages/rockchip-firmware/system.d/btattach@.service
[Unit]
Description=Attach /dev/ttyS0 to BlueZ stack using %I protocol
Wants=bluetooth.service
Before=bluetooth.service
After=dev-ttyS0.device

[Service]
Type=simple
ExecStartPre=/usr/sbin/rfkill unblock bluetooth
ExecStart=/usr/bin/btattach -B /dev/ttyS0 -P %I
ExecStopPost=/usr/sbin/rfkill block bluetooth
