# projects/OdroidXU3/packages/xu-fan-control/system.d/fan-control.service
[Unit]
Description=Fan control

[Service]
ExecStart=/usr/bin/fan-control.sh

[Install]
WantedBy=retroarch.target
