# projects/Ayn/devices/Odin/packages/usb-gadget-scripts/system.d/umtp-responder.service
[Unit]
Description=Start MTP Service
Requires=usb-gadget-init.service
After=usb-gadget-init.service

[Service]
ExecStart=umtprd
Restart=always
RestartSec=0
                      
[Install]
WantedBy=multi-user.target
