# projects/Ayn/devices/Odin/packages/usb-gadget-scripts/system.d/usb-gadget-init.service
[Unit]
Description=Init Nintendo Switch Gadget Framework
Before=usb-gadget-finalize.service

[Service]
ExecStart=/usr/bin/usb-gadget.sh init

[Install]
WantedBy=multi-user.target
