# packages/lakka/lakka_tools/gpicase_safeshutdown/system.d/gpicase2-safeshutdown.service
[Unit]
Description=GPICase2 safe shutdown

[Service]
ExecStart=/usr/bin/safeshutdown_gpicase2.py
TimeoutStopSec=5
Restart=always
RestartSec=5
StartLimitInterval=0

[Install]
WantedBy=multi-user.target
