# packages/linux-firmware/qca-firmware/system.d/qca-firmware.service
[Unit]
Description=QCA Bluetooth Firmware Service
ConditionPathExists=/dev/ttyAML1

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=/usr/bin/hciattach -n -s 115200 /dev/ttyAML1 qca 2000000
Restart=always
RestartSec=2
StartLimitInterval=0
LimitNPROC=1
TimeoutStopSec=1s
