# packages/addons/service/usbmuxd/source/system.d/service.usbmuxd.service
[Unit]
Description=Socket daemon for the usbmux protocol used by Apple devices

[Service]
Type=simple
ExecStartPre=/sbin/modprobe ipheth
ExecStart=/storage/.kodi/addons/service.usbmuxd/bin/usbmuxd --systemd
PIDFile=/var/run/usbmuxd.pid

[Install]
WantedBy=multi-user.target

