# packages/sysutils/open-vm-tools/system.d/open-vm-tools.service
[Unit]
Description=OpenVMTools Server
After=xorg-server.service

ConditionPathExists=/usr/bin/vmware-checkvm
ConditionPathExists=/usr/bin/vmtoolsd

[Service]
Type=forking
PIDFile=/var/run/vmtoolsd.pid
ExecStart=-/bin/sh -c '/usr/bin/vmware-checkvm && /usr/bin/vmtoolsd --background /var/run/vmtoolsd.pid'
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=graphical.target
