# packages/x11/driver/xf86-video-virtualbox/system.d/virtualbox-display.service
[Unit]
Description=VirtualBox Client Display Daemon
ConditionVirtualization=oracle
After=systemd-modules-load.service

[Service]
Environment=DISPLAY=:0.0
ExecStart=/usr/bin/VBoxClient --display --nodaemon
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical.target
