# projects/Switch/packages/mergerfs/system.d/tmp-cores.mount
[Unit]
Description=Cores directory
Before=retroarch.target
After=storage.mount
After=systemd-tmpfiles-setup.service

[Mount]
What=/usr/lib/libretro:/storage/cores
Where=/tmp/cores
Type=mergerfs
Options=defaults,allow_other,use_ino

[Install]
WantedBy=retroarch.target
