# projects/Switch/packages/mergerfs/system.d/tmp-assets.mount
[Unit]
Description=Assets directory
Before=retroarch.target
After=storage.mount
After=systemd-tmpfiles-setup.service

[Mount]
What=/usr/share/retroarch-assets:/storage/assets
Where=/tmp/assets
Type=mergerfs
Options=defaults,allow_other,use_ino

[Install]
WantedBy=retroarch.target
