# packages/sysutils/timezone-data/system.d/tz-data-monitor.service
[Unit]
Description=Timezone data monitor

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart tz-data.service

[Install]
WantedBy=multi-user.target
