# packages/sysutils/timezone-data/system.d/tz-data-monitor.path
[Unit]
Description=Timezone data monitor
After=local-fs.target

[Path]
PathModified=/storage/.xbmc/userdata/guisettings.xml

[Install]
WantedBy=multi-user.target
