# packages/sysutils/timezone-data/system.d/tz-data.service
[Unit]
Description=Setup Timezone data
Before=xbmc.service

ConditionPathExists=/storage/.xbmc/userdata/guisettings.xml

[Service]
Type=oneshot
ExecStart=/usr/lib/openelec/tzdata-setup
RemainAfterExit=yes
StartLimitInterval=0

[Install]
WantedBy=xbmc.service
