# packages/sysutils/timezone-data/system.d/tz-data.service
[Unit]
Description=Setup Timezone data
Before=xbmc.service

[Service]
Type=oneshot
Environment=TIMEZONE=UTC
EnvironmentFile=-/storage/.cache/timezone
ExecStart=/bin/ln -sf /usr/share/zoneinfo/${TIMEZONE} /var/run/localtime
RemainAfterExit=yes
StartLimitInterval=0

[Install]
WantedBy=xbmc.service
