# packages/sysutils/amremote/system.d/amlogic-remotecfg.service
[Unit]
Description=Amlogic IR remote support

[Service]
Type=oneshot
ExecStart=/bin/sh -c '[ -f /storage/.config/remote.conf ] && /usr/bin/remotecfg /storage/.config/remote.conf || /usr/bin/remotecfg /etc/remote.conf'
RemainAfterExit=yes

[Install]
WantedBy=basic.target
