# packages/sysutils/amremote/system.d/amlogic-remotecfg.service
[Unit]
Description=Amlogic IR remote support

[Service]
Type=oneshot
ExecStart=/usr/lib/openelec/remote-config
RemainAfterExit=yes

[Install]
WantedBy=basic.target
