# packages/addons/service/mariadb/source/profile.d/99-mysql-histfile.profile
# SPDX-License-Identifier: GPL-2.0
# Copyright (C) 2018-present Team LibreELEC (https://libreelec.tv)

export MYSQL_HISTFILE=/storage/.kodi/userdata/addon_data/service.mariadb/mysql_history
