# packages/mediacenter/OpenELEC-settings/system.d/backup-restore.service
[Unit]
Description=Restoring Backup
DefaultDependencies=no

[Service]
Type=idle
ExecStart=/usr/lib/openelec/backup-restore
StandardInput=tty-force
StandardOutput=inherit
StandardError=inherit
