# packages/mediacenter/OpenELEC-settings/system.d/factory-reset.service
[Unit]
Description=Factory reset
DefaultDependencies=false
After=show-version.service
Before=local-fs-pre.target shutdown.target
Conflicts=shutdown.target

ConditionPathExists=|/storage/.cache/reset_oe
ConditionPathExists=|/storage/.cache/reset_xbmc

[Service]
Type=oneshot
ExecStart=/usr/lib/openelec/factory-reset
StandardOutput=tty

[Install]
WantedBy=local-fs-pre.target
