# packages/x11/xserver/xorg-server/system.d/xorg.target
[Unit]
Description=The basic Xorg server
Conflicts=installer.target textmode.target rescue.target
Before=graphical.target
