# packages/x11/xserver/xorg-server/system.d/xorg-configure@.service
[Unit]
Description=configure Xorg Server for %I

ConditionKernelCommandLine=!installer

[Service]
Type=oneshot
ExecStart=/usr/lib/xorg/xorg-configure %I
RemainAfterExit=yes
