# packages/x11/xserver/xorg-server/system.d/xorg-server.service
[Unit]
Description=Xorg Server
Conflicts=rescue.service rescue.target installer.target
Requisite=graphical.target
After=multi-user.target

# On OpenELEC X11 is on tty1. We explicitly cancel the getty here to
# avoid any races around that.
#Conflicts=getty@tty1.service
#After=getty@tty1.service

[Service]
Environment=DISPLAY=:0.0
EnvironmentFile=/run/X11/xorg-settings.conf
ExecStart=/usr/bin/Xorg vt01 $X11_ARGS
Restart=always

[Install]
WantedBy=graphical.target
