# packages/x11/xserver/xorg-server/system.d/xorg-monitor.path
[Unit]
Description=Xorg configfile monitor
After=xorg-server.service

[Path]
PathModified=/run/openelec/debug.conf

[Install]
WantedBy=graphical.target
