# packages/x11/xserver/xorg-server/system.d/xorg.service
[Unit]
Description=Xorg Server
Before=graphical.target kodi.service
After=multi-user.target
ConditionKernelCommandLine=!installer

[Service]
Type=notify
EnvironmentFile=-/run/openelec/debug/xorg.conf
ExecStart=/usr/bin/xorg-launch vt01 $XORG_DEBUG
Restart=always
RestartSec=2
StartLimitInterval=0

[Install]
WantedBy=graphical.target
