# packages/x11/xserver/xorg-server/system.d/xorg-server.path
[Unit]
Description=Xorg server configfile watcher

[Path]
PathExists=/run/openelec/xorg-settings.conf

[Install]
WantedBy=graphical.target
