# packages/x11/xserver/xorg-server/system.d/xorg-server.service
[Unit]
Description=Xorg Server
Conflicts=rescue.service rescue.target installer.target
Requires=graphical.target
After=multi-user.target

[Service]
Environment=DISPLAY=:0.0
EnvironmentFile=/run/X11/xorg-settings.conf
ExecStart=/usr/bin/Xorg vt01 $X11_ARGS
Restart=always

[Install]
WantedBy=graphical.target
