# packages/x11/xserver/xorg-server/system.d/xorg-configure@.service
[Unit]
Description=Xorg Server configuration
Before=graphical.target

[Service]
ExecStart=/lib/systemd/xorg-configure %I
RemainAfterExit=yes
Type=oneshot
