# packages/x11/xserver/xorg-server/system.d/xorg-server.service
[Unit]
Description=Xorg Server
Conflicts=rescue.service rescue.target installer.target
Requires=graphical.target
After=multi-user.target

[Service]
Environment=DISPLAY=:0.0
EnvironmentFile=/run/openelec/xorg-settings.conf
EnvironmentFile=-/run/openelec/debug.conf
ExecStart=/usr/bin/Xorg vt01 $XORG_ARGS $XORG_DEBUG
Restart=always

[Install]
WantedBy=graphical.target
