# packages/x11/xserver/xorg-server/system.d/xorg-monitor.service
[Unit]
Description=Xorg configfile monitor
After=xorg.service

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart xorg.service

[Install]
WantedBy=graphical.target
