# projects/S912/packages/rng-tools/system.d/rng-tools.service
[Unit]
Description=RNG-Tools service
ConditionPathExists=/dev/hwrng

[Service]
Type=simple
ExecStart=/usr/sbin/rngd -f

[Install]
WantedBy=basic.target
