# packages/addons/service/tvheadend43/source/system.d/service.tvheadend43.service
[Unit]
Description=TVHeadend43 Service
After=network-online.service

[Service]
ExecStart=/bin/sh -c "exec sh /storage/.kodi/addons/service.tvheadend43/bin/tvheadend43.start"
TimeoutStopSec=2
Restart=always
RestartSec=2
StartLimitInterval=0

[Install]
WantedBy=kodi.target
