# packages/security/polkit/system.d/polkit.service
[Unit]
Description=PolicyKit daemon
After=syslog.target

[Service]
ExecStart=/usr/lib/polkit-1/polkitd --no-debug
KillMode=process

[Install]
WantedBy=multi-user.target
