# projects/Switch/devices/Switch-L4T/filesystem/usr/lib/udev/rules.d/91-xorg-conf-tegra.rules
SUBSYSTEM!="soc", GOTO="soc_end"
KERNEL!="soc0", GOTO="soc_end"

ATTR{machine}=="jetson_e", RUN+="/bin/ln -sf /etc/X11/xorg.conf.jetson_e /etc/X11/xorg.conf"

LABEL="soc_end"
