# packages/addons/service/minidlna/source/system.d/service.minidlna.service
[Unit]
Description=minidlna
After=graphical.target

[Service]
ExecStart=/bin/sh -c "exec sh /storage/.kodi/addons/service.minidlna/bin/minidlna.start"
TimeoutStopSec=1
Restart=always
RestartSec=2
StartLimitInterval=0

[Install]
WantedBy=kodi.target
