# projects/L4T/packages/xorg-server/system.d/xorg-configure@.service
[Unit]
Description=configure Xorg Server for %i

ConditionKernelCommandLine=!installer

[Service]
Type=oneshot
ExecStart=/usr/lib/xorg/xorg-configure %i
RemainAfterExit=yes
