# packages/graphics/libmali/system.d/libmali-setup.service
[Unit]
Description=Setup symlink for ARM Mali library
Before=graphical.target
After=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/libmali-setup

[Install]
WantedBy=multi-user.target
