# packages/network/bluez/udev.d/09-bluetooth.rules
ACTION!="add|change", GOTO="end"

KERNEL=="hci[0-9]*", SUBSYSTEM=="bluetooth", RUN+="/usr/bin/hciconfig %k up"
KERNEL=="hci[0-9]*", SUBSYSTEM=="bluetooth", RUN+="/usr/bin/hciconfig %k lm master"

LABEL="end"
