# packages/addons/service/spotify-connect-web/source/system.d/service.spotify-connect-web.service
[Unit]
Description=spotify-connect-web
After=network-online.target
Requires=network-online.target

[Service]
ExecStart=/bin/sh /storage/.kodi/addons/service.spotify-connect-web/bin/spotify-connect-web.start
Restart=on-failure

[Install]
WantedBy=kodi.target
