# packages/addons/driver/sapphire/source/system.d/driver.remote.sapphire.service
[Unit]
Description=Sapphire remote start script
Before=kodi.service

[Service]
Type=oneshot
Environment=HOME=/storage
ExecStartPre=/storage/.kodi/addons/driver.remote.sapphire/bin/sapphire-config
ExecStart=/storage/.kodi/addons/driver.remote.sapphire/bin/sapphire_startup.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
