# packages/addons/service/multimedia/vdr-addon/source/system.d/service.multimedia.vdr-addon.service
[Unit]
Description=vdr
After=graphical.target

[Service]
ExecStart=/bin/sh -c "exec sh /storage/.xbmc/addons/service.multimedia.vdr-addon/bin/vdr.start"
TimeoutStopSec=2
Restart=always
RestartSec=2
StartLimitInterval=0

[Install]
WantedBy=xbmc.target
