# packages/sysutils/lcdproc/system.d/lcdd.service
[Unit]
Description=LCDProc
ConditionPathExists=|/storage/.cache/services/lcdd.conf
After=multi-user.target

[Service]
EnvironmentFile=-/storage/.cache/services/lcdd.conf
ExecStart=/bin/sh -c 'exec /usr/lib/openelec/lcd-wrapper'
TimeoutStopSec=1s

[Install]
WantedBy=multi-user.target
