# packages/sysutils/lcdproc/system.d/lcdd.service
[Unit]
Description=LCDProc
ConditionPathExists=|/storage/.cache/services/lcdd.conf
After=syslog.target

[Service]
EnvironmentFile=-/storage/.cache/services/lcdd.conf
ExecStart=/bin/sh -c 'LCDCONF=/etc/LCDd.conf; [ -f /storage/.config/LCDd.conf ] && LCDCONF=/storage/.config/LCDd.conf; exec /usr/sbin/LCDd -f -c $LCDCONF -d $LCD_DRIVER'
TimeoutStopSec=1s

[Install]
WantedBy=multi-user.target
