# packages/sysutils/lirc/system.d/lircd-uinput.service
[Unit]
Documentation=http://lirc.org/html/configure.html
Description=Forward LIRC button presses as uinput events
ConditionPathExists=/storage/.config/lircd.conf
After=lircd.service

[Service]
Type=simple
ExecStart=/usr/lib/libreelec/lircd_uinput_helper --add-release-events

[Install]
WantedBy=multi-user.target
