# packages/sysutils/lirc/system.d/lircd@.service
[Unit]
Description=Lirc with %I
After=lircd-defaults.service
Requires=lircd-defaults.service

ConditionPathExists=/storage/.cache/services/lircd.conf

[Service]
Type=forking
ExecStart=/usr/lib/libreelec/lircd_helper %I
TimeoutStopSec=1
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target

