# packages/sysutils/lirc/system.d/lircd-uinput@.service
[Unit]
Description=lircd-uinput with %I

ConditionPathExists=/storage/.cache/services/lircd.conf

[Service]
ExecStart=/usr/sbin/lircd-uinput /run/lirc/lircd-%I
Slice=system-lircd.slice
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target

