# packages/sysutils/lirc/system.d/lircd-defaults.service
[Unit]
Description=Lirc defaults
After=local-fs.target

ConditionPathExists=!/storage/.cache/services/lircd.conf
ConditionPathExists=!/storage/.cache/services/lircd.disabled

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'cp /usr/share/services/lircd.conf /storage/.cache/services/'
RemainAfterExit=yes
