# packages/sysutils/lirc/system.d/lircd.socket
[Socket]
ListenStream=/var/run/lirc/lircd.socket

[Install]
WantedBy=sockets.target
Also=lircd.service
