# packages/sysutils/lirc/system.d/lircd.service
[Unit]
Documentation=http://lirc.org/html/configure.html
Description=Flexible IR remote input/output application support
ConditionPathExists=/storage/.config/lircd.conf
After=network.target

[Service]
Type=simple
ExecStart=/usr/lib/libreelec/lircd_helper --nodaemon

[Install]
WantedBy=multi-user.target
