# packages/tools/installer/system.d/installer.service
[Unit]
Description=XBMC Media Center
Requires=installer.target

[Service]
ExecStart=/usr/bin/installer
Type=idle
StandardInput=tty-force
StandardOutput=inherit
StandardError=inherit
KillMode=process

# Bash ignores SIGTERM, so we send SIGHUP instead, to ensure that bash
# terminates cleanly.
KillSignal=SIGHUP

[Install]
WantedBy=installer.target
