# packages/tools/installer/system.d/installer.target
[Unit]
Description=OpenELEC installer
Requires=basic.target
After=basic.target
Conflicts=rescue.target multi-user.target graphical.target
AllowIsolate=yes
