# packages/libretro/retroarch/udev.d/99-bluetooth-power-on.rules
ACTION=="add", KERNEL=="hci0", RUN+="/usr/bin/hciconfig hci0 up"
