# packages/sysutils/atvclient/udev.d/98-appleir.rules
ACTION=="add", SUBSYSTEM=="usb", \
  ATTR{idVendor}=="05ac", ATTR{idProduct}=="824[012]", \
  RUN+="wakeup_enable", \
  TAG+="systemd", ENV{SYSTEMD_WANTS}+="atvclient.service"
