# packages/addons/service/dispmanx_vnc/source/system.d/service.system.dispmanx_vnc.service
[Unit]
Description=dispmanx_vnc
After=kodi.service

[Service]
ExecStart=/bin/sh /storage/.kodi/addons/service.system.dispmanx_vnc/bin/dispmanx_vncserver-service
TimeoutStopSec=1
Restart=always
RestartSec=2
StartLimitInterval=0

[Install]
WantedBy=kodi.target
