# packages/lakka/L4T/system.d/tmp-system.mount
[Unit]
Description=RetroArch system directory
Before=retroarch.service
After=systemd-tmpfiles-setup.service

[Mount]
What=/storage/system:/usr/share/retroarch-system
Where=/storage/.tmp/system-workdir
Type=mergerfs
Options=defaults,allow_other,use_ino

[Install]
WantedBy=retroarch.target
