# packages/lakka/L4T/system.d/tmp-assets.mount
[Unit]
Description=Assets directory
Before=retroarch.service
After=storage.mount
After=systemd-tmpfiles-setup.service

[Mount]
What=/storage/assets:/usr/share/retroarch-assets
Where=/tmp/assets
Type=mergerfs
Options=defaults,allow_other,use_ino

[Install]
WantedBy=retroarch.target
