# packages/lakka/L4T/system.d/var-bluetoothconfig.mount
[Unit]
Description=mount bluez filesystem
Before=bluetooth.service
After=storage.mount systemd-tmpfiles-setup.service
DefaultDependencies=no
Conflicts=umount.target shutdown.target

[Mount]
What=/storage/.cache/bluetooth
Where=/var/bluetoothconfig
Type=rewritefs
Options=allow_other,config=/etc/rewritefs.conf

[Install]
WantedBy=multi-user.target

