# packages/linux-firmware/brcmfmac_sdio-firmware-rpi/system.d/brcmfmac_sdio-firmware.service
[Unit]
Description=Broadcom sdio firmware update for BCM43430A1
ConditionPathExists=/dev/serial1
After=network.target

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=/usr/bin/rpi-btuart

[Install]
WantedBy=network.target
