# packages/linux-firmware/brcmfmac_sdio-firmware-rpi/system.d/brcmfmac_sdio-firmware.service
[Unit]
Description=Broadcom sdio firmware update for BCM43430A1
ConditionFileNotEmpty=/proc/device-tree/soc/gpio@7e200000/bt_pins/brcm,pins
Requires=dev-serial1.device
After=dev-serial1.device network.target

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=/usr/bin/rpi-btuart

[Install]
WantedBy=network.target
