# packages/linux-firmware/brcmfmac_sdio-firmware-rpi/system.d/brcmfmac_sdio-firmware.service
[Unit]
Description=Broadcom sdio firmware update for BCM43430A1
ConditionPathExists=/dev/ttyAMA0
ConditionPathExists=/proc/device-tree/soc/gpio@7e200000/bt_pins
After=network.target

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=/usr/bin/hciattach /dev/ttyAMA0 bcm43xx 921600 noflow -

[Install]
WantedBy=network.target
