# packages/lakka/lakka_tools/gpicase_change_audio_device/system.d/gpicase-change-audio_device.service
[Unit]
Description=Change audio_device parameter in /.config/retroarch/retroarch.cfg
Before=retroarch-autostart.service

[Service]
Type=oneshot
ExecStart=/usr/bin/gpicase-change-audio_device.py
Restart=no
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
