# projects/S8X2/packages/linux-firmware/rtl8723bs_bt/system.d/rtl8723bs_bt-firmware.service
[Unit]
Description=Realtek sdio firmware update
ConditionPathExists=/dev/ttyS2

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=/usr/bin/rtk_hciattach -n -s 115200 /dev/ttyS2 rtk_h5 2000000
