# packages/sysutils/udev/config/rules.d/10_firmware.rules
ACTION=="add",  SUBSYSTEM=="firmware",  ENV{FIRMWARE}=="*",   RUN+="/sbin/firmware_helper"
