# packages/network/openssh/system.d/sshd-keygen.service
[Unit]
Description=OpenSSH Server Key Generation
ConditionPathExists=|!/storage/.cache/ssh/ssh_host_key
ConditionPathExists=|!/storage/.cache/ssh/ssh_host_rsa_key
ConditionPathExists=|!/storage/.cache/ssh/ssh_host_dsa_key

[Service]
ExecStart=/usr/sbin/sshd-keygen
Type=oneshot

