# packages/network/openssh/system.d/sshd@.service
[Unit]
Description=OpenSSH per-connection server daemon
Wants=sshd-keygen.service
After=sshd-keygen.service

[Service]
# EnvironmentFile=-/etc/sysconfig/sshd
ExecStart=-/usr/sbin/sshd
StandardInput=socket

