# packages/network/openssh/system.d/sshd.service
[Unit]
Description=OpenSSH server daemon
After=syslog.target network.target

[Service]
ExecStartPre=/usr/sbin/sshd-keygen
ExecStart=/usr/sbin/sshd -D
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process

[Install]
WantedBy=multi-user.target
