# packages/mediacenter/LibreELEC-settings/system.d/backup-restore.target
[Unit]
Description=Restoring Backup target
Requires=backup-restore.service
After=backup-restore.service
AllowIsolate=yes
