# packages/mediacenter/LibreELEC-settings/system.d/factory-reset.service
[Unit]
Description=Factory reset
DefaultDependencies=no

[Service]
Type=idle
ExecStart=/usr/lib/libreelec/factory-reset
StandardInput=tty-force
StandardOutput=inherit
StandardError=inherit
