# packages/mediacenter/LibreELEC-settings/system.d/factory-reset.target
[Unit]
Description=Factory reset target
Requires=factory-reset.service
After=factory-reset.service
AllowIsolate=yes
