# packages/lakka/retroarch_base/retroarch/system.d/retroarch.target
[Unit]
Description=Retroarch
Requires=multi-user.target graphical.target
After=graphical.target
#Requires=basic.target
#After=basic.target
Conflicts=rescue.target
AllowIsolate=yes

[Install]
Alias=default.target
