# packages/lakka/retroarch_base/retroarch/system.d/retroarch-autostart.service
[Unit]
Description=RetroArch user autostart script
Before=retroarch.service
After=graphical.target

[Service]
Type=oneshot
Environment=HOME=/storage
ExecStart=-/bin/sh -c ". /etc/profile; exec /bin/sh /storage/.config/autostart.sh"
RemainAfterExit=yes

[Install]
WantedBy=retroarch.service
