# packages/lakka/retroarch_base/retroarch/system.d/tmp-cheats.mount
[Unit]
Description=Cheats directory
Before=retroarch.service
Conflicts=umount.target shutdown.target
After=systemd-tmpfiles-setup.service
DefaultDependencies=no

[Mount]
What=none
Where=/tmp/cheats
Type=overlay
Options=lowerdir=/usr/share/libretro-cheats,upperdir=/storage/cheats,workdir=/storage/.tmp/cheats-workdir

[Install]
WantedBy=retroarch.target
