# packages/lakka/retroarch_base/retroarch/system.d/tmp-assets.mount
[Unit]
Description=Assets directory
Before=retroarch.service
Conflicts=umount.target shutdown.target
After=systemd-tmpfiles-setup.service
DefaultDependencies=no

[Mount]
What=none
Where=/tmp/assets
Type=overlay
Options=lowerdir=/usr/share/retroarch/assets,upperdir=/storage/assets,workdir=/storage/.tmp/assets-workdir

[Install]
WantedBy=retroarch.target
