# packages/lakka/retroarch_base/retroarch/system.d/tmp-cores.mount
[Unit]
Description=Cores directory
Before=retroarch.service umount.target
Conflicts=umount.target
After=systemd-tmpfiles-setup.service
DefaultDependencies=no

[Mount]
What=none
Where=/tmp/cores
Type=overlay
Options=lowerdir=/usr/lib/libretro,upperdir=/storage/custom/cores,workdir=/storage/.tmp/cores-workdir

[Install]
WantedBy=retroarch.service
