# projects/L4T/packages/systemd/system.d/envconfig.service
[Unit]
Description=Generate system-wide /etc/environment file
DefaultDependencies=no
After=systemd-tmpfiles-setup.service

[Service]
Type=oneshot
ExecStart=/usr/bin/environment-setup
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
