# projects/L4T/packages/systemd/system.d/hwdb.service
[Unit]
Description=Update hwdb.bin
DefaultDependencies=no
Conflicts=shutdown.target
After=systemd-tmpfiles-setup.service
Before=sysinit.target shutdown.target systemd-udevd.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=systemd-hwdb update

[Install]
WantedBy=sysinit.target
