# projects/L4T/packages/systemd/system.d/systemd-timesyncd-setup.service
[Unit]
Description=Setup NTP servers for timesyncd
DefaultDependencies=no
After=systemd-remount-fs.service systemd-sysusers.service

[Service]
Type=oneshot
ExecStart=/usr/sbin/systemd-timesyncd-setup
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
