# projects/L4T/packages/systemd/system.d/usercache.service
[Unit]
Description=Setup User cache dir
DefaultDependencies=no
After=systemd-tmpfiles-setup.service

[Service]
Type=oneshot
ExecStart=/usr/bin/usercache-setup
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
