# projects/L4T/packages/systemd/system.d/userconfig.service
[Unit]
Description=Setup User config dir
DefaultDependencies=no
After=systemd-tmpfiles-setup.service

[Service]
Type=oneshot
ExecStart=/usr/bin/userconfig-setup
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
