# projects/Gamegirl/packages/gamegirl-joypad/system.d/gamegirl-joypad.service
[Unit]
Description=Gamegirl joypad
After=retroarch.target
Requires=retroarch.target

[Service]
ExecStart=/usr/bin/gamegirl-joypad

[Install]
WantedBy=retroarch.target
