# packages/lakka/RPi/system.d/disable-hdmi.service
[Unit]
Description=Disable Raspberry Pi HDMI port

[Service]
Type=oneshot
ExecStart=/usr/bin/tvservice -o
ExecStop=/usr/bin/tvservice -p
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
