# packages/linux-drivers/wetekdvb/system.d/wetekdvb.service
[Unit]
Description=WeTek DVB module loader
ConditionPathExists=/proc/device-tree/dvb/dev_name
After=kernel-overlays.service

[Service]
Type=simple
ExecStart=/bin/sh -c '[ `cat /proc/device-tree/dvb/dev_name` = "wetek-dvb" ] && /sbin/modprobe wetekdvb'

[Install]
WantedBy=basic.target
