# projects/S805/packages/linux-firmware/brcmfmac_sdio-firmware-aml/system.d/brcmfmac_sdio-firmware-aml.service
[Unit]
Description=Broadcom sdio firmware update
ConditionPathExists=/dev/ttyS2

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=/usr/bin/hciattach -n -s 115200 /dev/ttyS2 bcm43xx 2000000
