# projects/L4T/devices/Switch/packages/switch-bsp/system.d/alsa-init.service
[Unit]
Description=Restore alsa configs
Before=pulseaudio.service

[Service]
ExecStart=alsactl --file /storage/.config/asound.state restore
ExecStop=alsactl --file /storage/.config/asound.state store

[Install]
WantedBy=multi-user.target
