# projects/L4T/devices/Switch/packages/switch-bsp/system.d/var-bluetoothconfig.mount
[Unit]
Description=mount bluez filesystem
Before=bluetooth.service

[Mount]
What=/storage/.cache/bluetooth
Where=/var/bluetoothconfig
Type=rewritefs
Options=allow_other,config=/etc/rewritefs.conf

[Install]
WantedBy=multi-user.target

