# projects/L4T/devices/Switch/packages/switch-bsp/system.d/xorg-configure-switch.service
[Unit]
Description=configure Xorg Server for Nintendo Switch

ConditionKernelCommandLine=!installer

[Service]
Type=oneshot
ExecStart=/usr/lib/xorg/xorg-configure switch
RemainAfterExit=yes

[Install]
WantedBy=graphical.target
