# packages/addons/service/hyperion/source/system.d/service.hyperion.service
[Unit]
Description=Hyperion service
After=graphical.target

[Service]
Type=forking
ExecStart=/bin/sh -c "exec sh /storage/.kodi/addons/service.hyperion/bin/hyperiond.start"
ExecReload=/bin/kill -HUP $MAINPID
PIDFile=/run/hyperiond.pid
TimeoutStopSec=2
Restart=always
RestartSec=10
StartLimitInterval=0

[Install]
WantedBy=default.target
