# packages/sysutils/systemd/system.d/cpufreq.service
[Unit]
Description=Set CPU frequency governor and its tunables
ConditionPathExists=/usr/bin/cpufreq

[Service]
Type=oneshot
ExecStart=/usr/bin/cpufreq
RemainAfterExit=yes

[Install]
WantedBy=basic.target
