# packages/sysutils/systemd/system.d/debugconfig.service
[Unit]
Description=Setup debug config
DefaultDependencies=no
Before=local-fs.target shutdown.target
Conflicts=shutdown.target

ConditionKernelCommandLine=|debugging
ConditionPathExists=|/storage/.cache/debug.openelec

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'mkdir -p /run/openelec/debug; cp /usr/share/debugconf/*.conf /run/openelec/debug'

[Install]
WantedBy=local-fs.target
