# packages/sysutils/systemd/system.d/debugconfig.service
[Unit]
Description=Setup debug config
DefaultDependencies=no
Before=sysinit.target

ConditionKernelCommandLine=debugging

[Service]
Type=oneshot
ExecStart=/lib/systemd/openelec-debugconfig

[Install]
WantedBy=sysinit.target
