# packages/sysutils/systemd/system.d/kernel-overlays.service
[Unit]
Description=Create kernel modules tree from overlays
DefaultDependencies=no
After=var.mount
Before=systemd-udevd.service systemd-udev-trigger.service systemd-modules-load.service systemd-remount-fs.service swap.service machine-id.service kmod-static-nodes.service

[Service]
Type=oneshot
ExecStart=/usr/sbin/kernel-overlays-setup

[Install]
WantedBy=local-fs.target
