# packages/sysutils/systemd/system.d/machine-id.service
[Unit]
Description=Setup machine-id
DefaultDependencies=no
After=local-fs.target
Before=systemd-journald.service

[Service]
Type=oneshot
ExecStart=/usr/bin/systemd-machine-id-setup
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
