# packages/sysutils/systemd/system.d/machine-id.service
[Unit]
Description=Setup machine-id
DefaultDependencies=no
Conflicts=shutdown.target
Before=systemd-journald.service systemd-tmpfiles-setup-dev.service shutdown.target
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/bin/systemd-machine-id-setup
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
