# packages/sysutils/systemd/system.d/userconfig.service
[Unit]
Description=Setup User config dir
DefaultDependencies=no
Before=local-fs.target shutdown.target
Conflicts=shutdown.target

[Service]
Type=oneshot
ExecStart=/bin/sh -c "mkdir -p /storage/.config; false | cp -iR /usr/config/* /storage/.config/ &>/dev/null"
RemainAfterExit=yes

[Install]
WantedBy=local-fs.target
