# packages/sysutils/systemd/system.d/userconfig.service
[Unit]
Description=Setup User config dir
DefaultDependencies=no
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/lib/systemd/openelec-userconfig
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
