# packages/sysutils/systemd/system.d/userconfig.service
[Unit]
Description=Setup User config dir
DefaultDependencies=no
After=systemd-tmpfiles-setup.service

[Service]
Type=oneshot
ExecStart=/bin/sh -c "false | cp -iR /usr/config/* /storage/.config/ &>/dev/null"
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
