# packages/network/bluez/system.d/bluetooth.service
[Unit]
Description=Bluetooth service
Requires=bluetooth.target

ConditionPathExists=/storage/.cache/services/bluez.conf

[Service]
NotifyAccess=main
# EnvironmentFile=/storage/.cache/services/avahi.conf
EnvironmentFile=-/run/openelec/debug.conf
ExecStart=/usr/lib/bluetooth/bluetoothd $BLUETOOTH_DEBUG
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_BIND_SERVICE CAP_NET_RAW
LimitNPROC=1

[Install]
WantedBy=bluetooth.target
Alias=dbus-org.bluez.service
