# packages/network/bluez/system.d/bluetooth-monitor.path
[Unit]
Description=Bluetooth configfile monitor

[Path]
PathModified=/storage/.cache/services/bluez.conf
PathModified=/run/openelec/debug/bluez.conf

[Install]
WantedBy=bluetooth.target
