# packages/network/bluez/system.d/bluetooth-monitor.service
[Unit]
Description=Bluetooth configfile monitor
Requires=bluetooth.target

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart bluetooth.service

[Install]
WantedBy=bluetooth.target
