# packages/network/bluez/system.d/obex.service
[Unit]
Description=Bluetooth service
Requires=bluetooth.target bluetooth.service
After=bluetooth.service

ConditionPathExists=/storage/.cache/services/obex.conf
ConditionPathExists=/storage/.cache/services/bluez.conf

[Service]
EnvironmentFile=/storage/.cache/services/obex.conf
EnvironmentFile=-/run/openelec/debug.conf
ExecStart=/usr/lib/bluetooth/obexd $OBEX_ARGS $OBEX_DEBUG

[Install]
WantedBy=bluetooth.target
Alias=dbus-org.bluez.obex.service
