# packages/network/bluez/system.d/obex-monitor.service
[Unit]
Description=Bluetooth configfile monitor
Requires=bluetooth.target bluetooth.service

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart obex.service

[Install]
WantedBy=bluetooth.target
