# packages/network/bluez/system.d/obex-monitor.path
[Unit]
Description=OBEX configfile monitor

[Path]
PathModified=/storage/.cache/services/obexd.conf
PathModified=/run/openelec/debug.conf

[Install]
WantedBy=bluetooth.target
