# packages/network/bluez/system.d/obex.service
[Unit]
Description=Bluetooth service
Requires=bluetooth.target bluetooth.service
After=bluetooth.service

ConditionPathExists=/storage/.cache/services/obexd.conf
ConditionPathExists=/storage/.cache/services/bluez.conf

[Service]
EnvironmentFile=/storage/.cache/services/obexd.conf
EnvironmentFile=-/run/openelec/debug/obexd.conf
ExecStartPre=/bin/sh -c 'mkdir -p $OBEXD_ROOT'
ExecStart=/usr/lib/bluetooth/obexd -r $OBEXD_ROOT $OBEXD_DEBUG
TimeoutStopSec=1s

[Install]
WantedBy=bluetooth.target
