# packages/addons/service/emby/source/system.d/service.emby.service
[Unit]
Description=Emby - streaming server
Documentation=http://emby.media
After=network.target
Requires=network.target

[Service]
ExecStart=/bin/sh /storage/.kodi/addons/service.emby/bin/emby.start
Restart=always

[Install]
WantedBy=kodi.target

