# projects/Rockchip/packages/rtl8723b_bt/system.d/rtl8723b_bt.service
[Unit]
Description=Load rtl8723b_bt firmware on /dev/ttyS0
Wants=bluetooth.service
Before=bluetooth.service
After=dev-ttyS0.device

[Service]
Type=simple
ExecStartPre=/usr/sbin/rfkill unblock bluetooth
ExecStart=/usr/bin/rtk_hciattach -n -s 115200 ttyS0 rtk_h5
ExecStopPost=/usr/sbin/rfkill block bluetooth

[Install]
WantedBy=multi-user.target
