# packages/addons/service/snapclient/source/system.d/service.snapclient.service
[Unit]
Description=Snapclient
After=kodi.service network-online.target sound.target
Requires=kodi.service network-online.target sound.target

[Service]
Type=forking
ExecStart=/bin/sh /storage/.kodi/addons/service.snapclient/bin/snapclient.start

[Install]
WantedBy=kodi.target
