# packages/addons/service/snapserver/source/system.d/service.snapserver.service
[Unit]
Description=Snapserver
After=kodi.service network-online.target sound.target
Requires=kodi.service network-online.target sound.target

[Service]
Type=forking
ExecStart=/bin/sh /storage/.kodi/addons/service.snapserver/bin/snapserver.start

[Install]
WantedBy=kodi.target
