# packages/addons/service/minisatip/source/system.d/service.minisatip.service
[Unit]
Description=minisatip
After=network-online.service

[Service]
ExecStart=/bin/sh -c "exec sh /storage/.kodi/addons/service.minisatip/bin/minisatip.start"
TimeoutStopSec=2
Restart=always
RestartSec=2
StartLimitInterval=0

[Install]
WantedBy=kodi.target
