# projects/L4T/devices/Switch/packages/usb-gadget-scripts/system.d/umtp-responder.service
[Unit]
Description=Start MTP Service
Requires=usb-gadget-setup.service
After=usb-gadget-setup.service

[Service]
ExecStart=umtprd
Restart=always
RestartSec=0
                      
[Install]
WantedBy=multi-user.target
