# projects/L4T/devices/Switch/packages/usb-gadget-scripts/system.d/usb-gadget.service
[Unit]
Description=Configure Nintendo Switch Gadget Framework

[Service]
ExecStart=/usr/bin/usb-gadget.sh

[Install]
WantedBy=multi-user.target
