# projects/L4T/devices/Switch/packages/usb-gadget-scripts/system.d/usb-tty.service
[Unit]
Description=Start getty on ttyGS0
Requires=usb-gadget-setup.service
After=usb-gadget-setup.service

[Service]
ExecStart=agetty --local-line 9600 ttyGS0 vt100
Restart=always
RestartSec=0
                      
[Install]
WantedBy=multi-user.target
