# projects/L4T/devices/Switch/packages/usb-gadget-scripts/system.d/usb-gadget-finalize.service
[Unit]
Description=Finalize Nintendo Switch Gadget Framework
Requires=usb-gadget-init.service
Requires=umtp-responder.service
After=usb-gadget-init.service
After=umtp-responder.service

[Service]
ExecStart=/usr/bin/usb-gadget.sh finalize

[Install]
WantedBy=multi-user.target
