# projects/L4T/devices/Switch/packages/usb-gadget-scripts/system.d/usb-gadget-setup.service
[Unit]
Description=Setup Nintendo Switch Gadget Framework

[Service]
ExecStart=/usr/bin/usb-gadget.sh init

[Install]
WantedBy=multi-user.target
