# projects/L4T/devices/Switch/packages/usb-gadget-scripts/system.d/usb-tty.service
[Unit]
Description=Configure Nintendo Switch Gadget Framework
After=usb-gadget.service

[Service]
ExecStart=agetty --local-line 9600 ttyGS0 vt100
Restart=always
RestartSec=0
                      
[Install]
WantedBy=multi-user.target
