# packages/network/rpcbind/system.d/rpcbind.service
[Unit]
Description=RPC bind service
Requires=rpcbind.socket
After=syslog.target network.target

[Service]
Type=forking
# EnvironmentFile=/etc/sysconfig/rpcbind
ExecStart=/usr/bin/rpcbind -w ${RPCBIND_ARGS}

[Install]
WantedBy=multi-user.target
Also=rpcbind.socket

