# packages/sysutils/tz/system.d/tz-data.service
[Unit]
Description=Setup Timezone data
DefaultDependencies=no
Before=systemd-udevd.service
After=var.mount

[Service]
Type=oneshot
Environment=TIMEZONE=UTC
EnvironmentFile=-/storage/.cache/timezone
ExecStartPre=/bin/mkdir -p /var/run
ExecStart=/bin/ln -sf /usr/share/zoneinfo/${TIMEZONE} /var/run/localtime
RemainAfterExit=yes
StartLimitInterval=0

[Install]
WantedBy=sysinit.target
