# packages/addons/service/librespot/source/system.d/service.librespot.service
[Unit]
Description=librespot
After=network-online.target sound.target
Requires=network-online.target sound.target

[Service]
ExecStart=/bin/sh /storage/.kodi/addons/service.librespot/bin/librespot.start
ExecStopPost=/usr/bin/pactl suspend-sink librespot_sink 1
Restart=on-failure

[Install]
WantedBy=kodi.target
