# packages/addons/service/librespot/source/system.d/service.librespot.service
[Unit]
Description=librespot
Wants=kodi.service
After=kodi.service

[Service]
EnvironmentFile=/storage/.kodi/addons/service.librespot/bin/librespot.env
ExecStart=/bin/sh /storage/.kodi/addons/service.librespot/bin/librespot.start
ExecStopPost=/usr/bin/pactl suspend-sink "$LS_SINK" 1
Restart=on-failure

[Install]
WantedBy=kodi.service
