# packages/addons/service/librespot/source/system.d/service.librespot.service
[Unit]
Description=librespot
After=network-online.target
Requires=network-online.target

[Service]
ExecStart=/bin/sh /storage/.kodi/addons/service.librespot/bin/librespot.start
Restart=on-failure

[Install]
WantedBy=kodi.target
