# packages/addons/service/librespot/source/system.d/service.librespot.service
[Unit]
Description=librespot
Wants=kodi.service
After=kodi.service

[Service]
EnvironmentFile=/storage/.kodi/addons/service.librespot/ls_env.py
ExecStart=/bin/sh /storage/.kodi/addons/service.librespot/bin/librespot.start
ExecStopPost=-/usr/bin/pactl suspend-sink librespot_sink 1
Restart=always

[Install]
WantedBy=kodi.service
