# packages/addons/service/librespot/source/system.d/service.librespot.service
[Unit]
Description=librespot
<<<<<<< HEAD
After=network-online.target
Requires=network-online.target

[Service]
ExecStart=/bin/sh /storage/.kodi/addons/service.librespot/bin/librespot.start
=======
After=kodi.service network-online.target sound.target
Requires=kodi.service network-online.target sound.target

[Service]
EnvironmentFile=/storage/.kodi/addons/service.librespot/bin/librespot.env
ExecStart=/bin/sh /storage/.kodi/addons/service.librespot/bin/librespot.start
ExecStopPost=/usr/bin/pactl suspend-sink "$LS_SINK" 1
>>>>>>> 6ecb5bc094f11d18cb0b884bcf486a15b9a676d5
Restart=on-failure

[Install]
WantedBy=kodi.target
