# projects/L4T/devices/Switch/packages/retroarch/system.d/tmp-assets.mount
[Unit]
Description=Assets directory
Before=retroarch.service umount.target
Conflicts=umount.target
After=systemd-tmpfiles-setup.service
DefaultDependencies=no

[Mount]
What=/storage/assets:/usr/share/retroarch/assets
Where=/tmp/assets
Type=mergerfs
Options=defaults,allow_other,use_ino

[Install]
WantedBy=retroarch.target
