# projects/L4T/devices/Switch/packages/retroarch/system.d/tmp-overlays.mount
[Unit]
Description=Overlays directory
Before=retroarch.service umount.target
Conflicts=umount.target
After=systemd-tmpfiles-setup.service
DefaultDependencies=no

[Mount]
What=/storage/overlays:/usr/share/retroarch/overlays
Where=/tmp/overlays
Type=mergerfs
Options=defaults,allow_other,use_ino

[Install]
WantedBy=retroarch.target
