# projects/L4T/devices/Switch/packages/retroarch/system.d/tmp-database.mount
[Unit]
Description=Database directory
Before=retroarch.service umount.target
Conflicts=umount.target
After=systemd-tmpfiles-setup.service
DefaultDependencies=no

[Mount]
What=/storage/database:/usr/share/libretro-database
Where=/tmp/database
Type=mergerfs
Options=defaults,allow_other,use_ino

[Install]
WantedBy=retroarch.target
