# packages/addons/service/multimedia/tvheadend/source/system.d/service.multimedia.tvheadend.service
[Unit]
Description=TVHeadend Service
After=graphical.target

[Service]
ExecStart=/bin/sh -c "exec sh /storage/.kodi/addons/service.multimedia.tvheadend/bin/tvheadend.start"
TimeoutStopSec=2
Restart=always
RestartSec=2
StartLimitInterval=0

[Install]
WantedBy=kodi.target
