# packages/addons/service/tvmosaic/source/system.d/service.tvmosaic.service
[Unit]
Description=TV Mosaic Server
Documentation=https://tv-mosaic.com
After=network-online.target
Requires=network-online.target

[Service]
ExecStart=/bin/sh /storage/.kodi/addons/service.tvmosaic/bin/tvmosaic.start
Type=forking
Restart=always

[Install]
WantedBy=kodi.target
