# packages/sysutils/ethmactool/system.d/ethmactool-config.service
[Unit]
Description=Configure eth0 MAC address
BindsTo=sys-subsystem-net-devices-eth0.device
After=sys-subsystem-net-devices-eth0.device
Before=connman.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/ethmactool-config eth0

[Install]
WantedBy=multi-user.target
