# packages/audio/pulseaudio/system.d/pulseaudio.service
[Unit]
Description=PulseAudio Sound System
After=syslog.target local-fs.target

[Service]
ExecStart=/usr/bin/pulseaudio --system
Restart=always

[Install]
WantedBy=multi-user.target
