# packages/audio/pulseaudio/system.d/pulseaudio.service
[Unit]
Description=PulseAudio Sound Service
After=syslog.target local-fs.target

[Service]
Type=dbus
BusName=org.pulseaudio.Server
ExecStart=/usr/bin/pulseaudio --system
Restart=on-failure

[Install]
WantedBy=multi-user.target
