# packages/audio/pulseaudio/system.d/pulseaudio.service
[Unit]
Description=PulseAudio Sound System
After=syslog.target local-fs.target

[Service]
Restart=on-failure
ExecStart=/usr/bin/pulseaudio --system
TimeoutStopSec=1s
RestartSec=2
StartLimitInterval=60
StartLimitBurst=10

[Install]
WantedBy=multi-user.target
