# packages/audio/pulseaudio/system.d/pulseaudio.service
[Unit]
Description=PulseAudio Sound Service
After=syslog.target local-fs.target

[Service]
Type=dbus
BusName=org.pulseaudio.Server
ExecStart=/usr/bin/pulseaudio --system --disallow-exit --exit-idle-time=-1 --disable-shm --enable-memfd --realtime
Restart=on-failure
Environment=PULSE_STATE_PATH=/storage/.config/pulse
Environment=PULSE_CONFIG_PATH=/storage/.config/pulse

[Install]
WantedBy=multi-user.target
