# packages/linux-drivers/amlogic/fd628-aml/system.d/fd628.service
[Unit]
Description=Amlogic FD628 Service
ConditionPathExists=/proc/device-tree/fd628_dev

[Service]
Type=oneshot
ExecStart=/usr/sbin/modprobe aml_fd628
ExecStart=/usr/sbin/FD628Service
RemainAfterExit=yes

[Install]
WantedBy=basic.target
