# packages/audio/rpi-cirrus-config/udev.d/90-alsa-restore.rules
# SPDX-License-Identifier: GPL-2.0
# Copyright (C) 2016-present Team LibreELEC (https://libreelec.tv)

# When a sound device is detected, restore the volume settings
SUBSYSTEM=="sound", KERNEL=="controlC*", NAME="snd/%k", ACTION=="add", GOTO="alsa_restore_go"
GOTO="alsa_restore_end"

LABEL="alsa_restore_go"

# Separate config-script for RPi-Cirrus card
DRIVERS=="snd-rpi-cirrus", RUN+="rpi-cirrus-config", GOTO="alsa_restore_end"

# Default config-script for all other cards
RUN+="soundconfig %k"

LABEL="alsa_restore_end"

