# packages/sysutils/busybox/system.d/fs-resize.service
[Unit]
Description=FS Resize
DefaultDependencies=no

[Service]
Type=idle
ExecStart=/usr/lib/openelec/fs-resize
StandardInput=tty-force
StandardOutput=inherit
StandardError=inherit
