# packages/sysutils/busybox/system.d/fs-resize.service
[Unit]
Description=FS Resize
DefaultDependencies=false
After=show-version.service
Before=local-fs-pre.target shutdown.target
Conflicts=shutdown.target

ConditionPathExists=|/storage/.please_resize_me

[Service]
Type=oneshot
ExecStart=/usr/lib/openelec/fs-resize
StandardOutput=tty

[Install]
WantedBy=local-fs-pre.target
