# packages/sysutils/busybox/system.d/ledfix.service
[Unit]
Description=LEDfix Service
After=network-online.target graphical.target
Requires=graphical.target
Wants=kodi.target

[Service]
Type=oneshot
ExecStart=/bin/sh /usr/bin/ledfix
RemainAfterExit=yes
StartLimitInterval=0

[Install]
WantedBy=kodi.target
