# packages/sysutils/busybox/system.d/rpi-flash-firmware.service
[Unit]
Description=RPi Flash Firmware
Requires=tmp.mount var.mount sys-kernel-config.mount kernel-overlays.service
After=tmp.mount var.mount sys-kernel-config.mount kernel-overlays.service
DefaultDependencies=no

[Service]
Type=idle
ExecStart=/usr/lib/libreelec/rpi-flash-firmware
StandardInput=tty-force
StandardOutput=inherit
StandardError=inherit
