# packages/sysutils/busybox/system.d/rpi-flash-firmware.target
[Unit]
Description=RPi Flash Firmware target
Requires=rpi-flash-firmware.service
After=rpi-flash-firmware.service
AllowIsolate=yes
