# packages/sysutils/busybox/system.d/show-version.service
[Unit]
Description=Show Version
DefaultDependencies=no
Before=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/bin/lsb_release
RemainAfterExit=yes
StandardOutput=tty

[Install]
WantedBy=sysinit.target
