# packages/sysutils/busybox/system.d/storage-log.service
[Unit]
Description=Create Persistent Log Directory on /storage
DefaultDependencies=no
RequiresMountsFor=/storage

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/mkdir -p /storage/.cache/log/journal
