# packages/sysutils/busybox/system.d/var-log.mount
[Unit]
Description=Persistent Log Storage
RequiresMountsFor=/var /storage
Requires=storage-log.service
After=storage-log.service
ConditionKernelCommandLine=!installer
ConditionKernelCommandLine=|debugging
ConditionPathExists=|/storage/.cache/debug.libreelec

[Mount]
What=/storage/.cache/log
Where=/var/log
Options=bind
LazyUnmount=yes
