# packages/sysutils/busybox/system.d/var-log-debug.service
[Unit]
Description=Debug /var/log relink
ConditionKernelCommandLine=debugging
ConditionKernelCommandLine=!installer

[Service]
Environment=HOME=/storage
ExecStart=/bin/sh -c 'rm -rf /var/log && mkdir -p /storage/log && ln -sf /storage/log /var/log'

[Install]
WantedBy=local-fs.target
