# packages/addons/service/touchscreen/source/system.d/service.touchscreen.service
[Unit]
Description=Touchscreen support daemon
After=graphical.target
Before=kodi.target

[Service]
Type=oneshot
Environment=HOME=/storage
ExecStartPre=-/bin/sh -c "exec /bin/sh /storage/.kodi/addons/service.touchscreen/bin/ts_calibrate.sh service"
ExecStart=-/bin/sh -c    "exec /bin/sh /storage/.kodi/addons/service.touchscreen/bin/ts_uinput_touch.sh service"
RemainAfterExit=yes

[Install]
WantedBy=kodi.target
