# packages/x11/other/ratpoison/system.d/ratpoison.service
[Unit]
Description=Ratpoison Window Manager
After=xorg-server.service
Requires=graphical.target

[Service]
Environment=DISPLAY=:0.0
ExecStart=/bin/sh -c 'sleep 1; exec /usr/bin/ratpoison'
KillMode=process
Restart=always

[Install]
WantedBy=graphical.target
