# packages/x11/other/ratpoison/system.d/ratpoison.service
[Unit]
Description=Ratpoison Window Manager
After=xorg.target
Requires=xorg.target graphical.target

[Service]
Environment=DISPLAY=:0
# ExecStart=/bin/sh -c 'sleep 1; exec /usr/bin/ratpoison'
ExecStart=/usr/bin/ratpoison --file /etc/ratpoisonrc
KillMode=process
Restart=always

[Install]
WantedBy=graphical.target
