# packages/x11/other/ratpoison/system.d/ratpoison.service
[Unit]
Description=Ratpoison Window Manager
Requires=xorg.service
After=xorg.service

[Service]
Environment=DISPLAY=:0
ExecStart=/usr/bin/ratpoison --file /etc/ratpoisonrc
KillMode=process
Restart=always
RestartSec=1

[Install]
WantedBy=graphical.target
