# packages/sysutils/udevil/udev.d/95-udevil-mount.rules
# Mount blockdevices, /dev/sd* and partitions only
SUBSYSTEM!="block", KERNEL!="sd*", ENV{DEVTYPE}!="partition", GOTO="exit"

# Mount the device
ACTION=="add", RUN+="/usr/bin/udevil --mount /dev/%k"

# Clean up after removal
ACTION=="remove", RUN+="/usr/bin/udevil --umount /dev/%k"

# Exit
LABEL="exit"

