# packages/mediacenter/xbmc-master/system.d/xbmc-hacks.service
[Unit]
Description=XBMC hacks
Before=xbmc.service

[Service]
Type=oneshot
Environment=HOME=/storage
ExecStart=/usr/lib/xbmc/xbmc-hacks
RemainAfterExit=yes

[Install]
WantedBy=xbmc.service
