# packages/mediacenter/xbmc-master/system.d/xbmc-halt.service
[Unit]
Description=XBMC halt script
After=xbmc.service
Before=systemd-halt.service
DefaultDependencies=no

[Service]
Type=oneshot
Environment=HOME=/storage
ExecStart=-/bin/sh /storage/.config/shutdown.sh halt
RemainAfterExit=yes

[Install]
WantedBy=halt.target
