# projects/L4T/packages/busybox/system.d/fs-resize.target
[Unit]
Description=FS Resize target
Requires=fs-resize.service
After=fs-resize.service
AllowIsolate=yes
