# projects/L4T/packages/busybox/system.d/textmode.target
[Unit]
Description=Textmode
Requires=multi-user.target
After=multi-user.target
Conflicts=rescue.target
AllowIsolate=yes
