# packages/sysutils/entropy/system.d/add-random-at-shutdown.service
[Unit]
Description=Save random entropy at shutdown
DefaultDependencies=no
Before=systemd-poweroff.service systemd-reboot.service systemd-halt.service

[Service]
Type=oneshot
ExecStart=/usr/lib/entropy/add-random-at-shutdown
RemainAfterExit=yes

[Install]
WantedBy=poweroff.target reboot.target halt.target
