# packages/sysutils/entropy/system.d/add-entropy.service
[Unit]
Description=Add random entropy from file
DefaultDependencies=no
After=systemd-tmpfiles-setup.service
Before=systemd-udevd.service

[Service]
Type=oneshot
ExecStart=/usr/lib/entropy/add-entropy
RemainAfterExit=yes

[Install]
WantedBy=basic.target
