# packages/addons/service/prometheus-node-exporter/source/system.d/service.prometheus-node-exporter.service
[Unit]
Description=Prometheus Node Exporter service
After=network-online.target
Requires=network-online.target

[Service]
ExecStart=/storage/.kodi/addons/service.prometheus-node-exporter/bin/prometheus-node-exporter.start
Restart=always

[Install]
WantedBy=kodi.target
