# projects/L4T/devices/Switch/packages/switch-gpu-profile/system.d/switch-gpu-profile.service
[Unit]
Description=Switch GPU Profile
ConditionKernelCommandLine=!installer
Requires=xorg.service
After=xorg.service

[Service]
Environment=TERM=linux
Type=oneshot
ExecStart=/usr/bin/gpu-profile init
StandardError=syslog
RemainAfterExit=no

[Install]
WantedBy=graphical.target
